# dynu_updater/rootfs/etc/s6-overlay/s6-rc.d/init/finish
#!/command/with-contenv bashio
DYNU_CONF="/run/dynu-updater.conf"

if [[ -f "$DYNU_CONF" ]]; then
	USE_MQTT=$(jq -r '.mqtt.enabled' "$DYNU_CONF")
	MQTT_HOST=$(jq -r '.mqtt.host' "$DYNU_CONF")
	MQTT_PORT=$(jq -r '.mqtt.port' "$DYNU_CONF")
	MQTT_USER=$(jq -r '.mqtt.user' "$DYNU_CONF")
	MQTT_PASSWORD=$(jq -r '.mqtt.password' "$DYNU_CONF")
	readarray -t TLS_FLAGS < <(jq -r '.mqtt.tls_flags[]' <<< "$API_DATA")

	# Publish offline availability
	if [[ "$USE_MQTT" == "true" ]]; then
		bashio::log.info "Dynu Updater Stopping - Sending offline signal to MQTT"
		mosquitto_pub \
		  -h "$MQTT_HOST" \
		  -u "$MQTT_USER" \
		  -p "$MQTT_PORT" \
		  -P "$MQTT_PASSWORD" \
		  -t "dynu_updater/availability" \
		  -m "offline" \
		  -r \
		  "${TLS_FLAGS[@]}"
	fi
fi
bashio::log.info "Service init exited"