# dynu_updater/rootfs/etc/s6-overlay/s6-rc.d/init/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Init
# ==============================================================================
bashio::log.info "Dynu Updater Initialisation - Generating service configuration"

if bashio::config.has_value 'dynu_hostname'; then
	DYNU_HOSTNAME=$(bashio::config 'dynu_hostname' | tr '[:upper:]' '[:lower:]' | tr -d '[:space:]')
else
	DYNU_HOSTNAME=""
fi
if bashio::config.has_value 'cert.email'; then
	EMAIL=$(bashio::config 'cert.email' | xargs)
else
	EMAIL=""
fi
if bashio::config.has_value 'cert.domains'; then
	CERT_DOMAINS=$(bashio::config 'cert.domains' | tr '[:upper:]' '[:lower:]' | xargs)
else
	CERT_DOMAINS=""
fi
if bashio::config.has_value 'dns.domains'; then
	IP_DOMAINS=$(bashio::config 'dns.domains' | tr '[:upper:]' '[:lower:]' | xargs)
else
	IP_DOMAINS=""
fi
if bashio::config.true 'verbose'; then
    VERBOSE="true"
else
    VERBOSE="false"
fi
if bashio::config.true 'cert.force_renew'; then
    FORCE_RENEW="true"
else
    FORCE_RENEW="false"
fi
if bashio::config.true 'cert.dry_run'; then
    DRY_RUN="true"
else
    DRY_RUN="false"
fi
if bashio::config.has_value 'dynu_api_key'; then
	DYNU_API_KEY=$(bashio::config 'dynu_api_key' | xargs)
else
	DYNU_API_KEY=""
fi
if bashio::config.has_value 'dns.dynu_update_password'; then
	PASSWORD=$(bashio::config 'dns.dynu_update_password' | xargs)
else
	PASSWORD=""
fi
if bashio::config.has_value 'cert.keyfile'; then
	KEYFILE=$(bashio::config 'cert.keyfile' | xargs)
else
	KEYFILE=""
fi
if bashio::config.has_value 'cert.certfile'; then
	CERTFILE=$(bashio::config 'cert.certfile' | xargs)
else
	CERTFILE=""
fi
if bashio::config.has_value 'cert.cafile'; then
	CAFILE=$(bashio::config 'cert.cafile' | xargs)
else
	CAFILE=""
fi
if bashio::config.has_value 'mqtt.host'; then
	MQTT_HOST=$(bashio::config 'mqtt.host' | xargs)
else
	MQTT_HOST=""
fi
if bashio::config.has_value 'mqtt.port'; then
	MQTT_PORT=$(bashio::config 'mqtt.port')
else
	MQTT_PORT=""
fi
if bashio::config.has_value 'mqtt.user'; then
	MQTT_USER=$(bashio::config 'mqtt.user' | xargs)
else
	MQTT_USER=""
fi
if bashio::config.has_value 'mqtt.password'; then
	MQTT_PASSWORD=$(bashio::config 'mqtt.password' | xargs)
else
	MQTT_PASSWORD=""
fi
if bashio::config.has_value 'mqtt.keyfile'; then
	MQTT_KEYFILE=$(bashio::config 'mqtt.keyfile' | xargs)
else
	MQTT_KEYFILE=""
fi
if bashio::config.has_value 'mqtt.certfile'; then
	MQTT_CERTFILE=$(bashio::config 'mqtt.certfile' | xargs)
else
	MQTT_CERTFILE=""
fi
if bashio::config.has_value 'mqtt.cafile'; then
	MQTT_CAFILE=$(bashio::config 'mqtt.cafile' | xargs)
else
	MQTT_CAFILE=""
fi
if bashio::config.has_value 'cert.update_cert_days'; then
	CERT_TIME=$(bashio::config 'cert.update_cert_days')
else
	CERT_TIME=5
fi
if bashio::config.true 'cert.update_cert'; then
    CERT_VALID="true"
else
    CERT_VALID="false"
fi
if bashio::config.true 'dns.update_ip'; then
    DNS_VALID="true"
else
    DNS_VALID="false"
fi
if bashio::config.true 'mqtt.enabled'; then
    USE_MQTT="true"
else
    USE_MQTT="false"
fi
if bashio::config.true 'mqtt.core'; then
    MQTT_CORE="true"
else
    MQTT_CORE="false"
fi
if bashio::config.true 'mqtt.tls'; then
    MQTT_USE_TLS="true"
else
    MQTT_USE_TLS="false"
fi
if bashio::config.true 'mqtt.server_cert'; then
    MQTT_VALIDATE="true"
else
    MQTT_VALIDATE="false"
fi

if [[ "$CERT_VALID" == "false" && "$DNS_VALID" == "false" ]]; then
	bashio::log.error "Dynu Updater Initialisation - Options to Update IP Addresses and Generate Certificate are both disabled - Nothing to do"
	exit 1
fi

if [[ -z "$DYNU_HOSTNAME" || ! "$DYNU_HOSTNAME" =~ ^[A-Za-z0-9.-]+\.[A-Za-z]{2,}(\.[A-Za-z]{2,})*$ ]]; then
    bashio::log.error "Dynu Updater Initialisation - Invalid Dynu Hostname"
	exit 1
fi

if [[ -z "$DYNU_API_KEY" ]]; then
    bashio::log.error "Dynu Updater Initialisation - Invalid Dynu API Key"
	exit 1
fi

if [[ "$DNS_VALID" == "true" ]]; then
	if [[ -z "$IP_DOMAINS" ]]; then
		bashio::log.error "Dynu Updater Initialisation - no DNS Update Hostnames configured"
		exit 1
	fi

	if [[ -z "$PASSWORD" ]]; then
		bashio::log.error "Dynu Updater Initialisation - no Dynu Account/IP Update Password configured"
		exit 1
	fi
fi

if [[ "$CERT_VALID" == "true" ]]; then
	if [[ -z "$EMAIL" || ! "$EMAIL" =~ ^[^@]+@[^@]+\.[^@]+$ ]]; then
		bashio::log.error "Dynu Updater Initialisation - invalid email address"
		exit 1
	fi

	if [[ -z "$CERT_DOMAINS" ]]; then
		bashio::log.error "Dynu Updater Initialisation - no Certificate Hostnames configured"
		exit 1
	fi

	if [[ -z "$CERTFILE" || ! "$CERTFILE" =~ ^[A-Za-z0-9._/-]+$ ]]; then
		bashio::log.error "Dynu Updater Initialisation - Invalid Certificate File"
		exit 1
	fi

	if [[ -z "$KEYFILE" || ! "$KEYFILE" =~ ^[A-Za-z0-9._/-]+$ ]]; then
		bashio::log.error "Dynu Updater Initialisation - Invalid Private Key File"
		exit 1
	fi
	
	if [[ -z "$CAFILE" || ! "$CAFILE" =~ ^[A-Za-z0-9._/-]+$ ]]; then
		bashio::log.error "Dynu Updater Initialisation - Invalid CA File"
		exit 1
	fi
	
	if [[ -z "$CERT_TIME" ]] || (( CERT_TIME < 1 || CERT_TIME > 30 )); then
		bashio::log.error "Dynu Updater Initialisation - Invalid Certificate Check Time"
		exit 1
	fi
fi

bashio::log.info "Dynu Updater Initialisation - Contacting Dynu API"

if DYNU_API_RESPONSE=$(curl -s --max-time 5 -H "API-Key: $DYNU_API_KEY" "https://api.dynu.com/v2/dns" 2>/dev/null); then
	HTTP_STATUS=$(jq -r '.statusCode // empty' <<< "$DYNU_API_RESPONSE")
	
	if [[ -z "$HTTP_STATUS" ]]; then
		bashio::log.error "Dynu Updater Initialisation - Invalid or empty JSON returned from Dynu Domain API"
		exit 1
	fi
	
	EXCEPTION_TYPE=$(jq -r '.exception.type // empty' <<< "$DYNU_API_RESPONSE")
	
	EXCEPTION_MESSAGE=$(jq -r '.exception.message // empty' <<< "$DYNU_API_RESPONSE")

	case "$HTTP_STATUS" in
		200)
			bashio::log.info "Dynu Updater Initialisation - Dynu Domain API data retrieved."
			;;
		*)
			bashio::log.error "Dynu Updater Initialisation - Dynu Domain API Error: status: $HTTP_STATUS, type: $EXCEPTION_TYPE, message: $EXCEPTION_MESSAGE"
			exit 1
			;;

	esac
else
	bashio::log.error "Dynu Updater Initialisation - Dynu Domain API could not be loaded"
	exit 1
fi

# Retrieve Domain ID
DYNU_DOMAIN_ID=$(
	jq -r --arg HOST "$DYNU_HOSTNAME" '
		.domains[]
		| select(.name == $HOST)
		| .id // empty
	' <<< "$DYNU_API_RESPONSE")

if [[ -z "$DYNU_DOMAIN_ID" ]]; then
	bashio::log.error "Dynu Updater Initialisation - '$DYNU_HOSTNAME' not found in the Dynu Domain API response"
	exit 1
fi

# Retrieve IPv6 enabled status for domain
DYNU_IPV6_ENABLED=$(jq -r --arg HOST "$DYNU_HOSTNAME" '
	.domains[]
	| select(.name == $HOST)
	| .ipv6 // "false"
' <<<"$DYNU_API_RESPONSE")

if [[ "$DYNU_IPV6_ENABLED" == "true" ]]; then
	bashio::log.info "Dynu Updater Initialisation - IPv6 enabled for $DYNU_HOSTNAME."
fi

if DYNU_RECORDS=$(curl -s --max-time 5 -H "API-Key: $DYNU_API_KEY" "https://api.dynu.com/v2/dns/${DYNU_DOMAIN_ID}/record" 2>/dev/null); then
	HTTP_STATUS=$(
	jq -r '.statusCode // empty' <<< "$DYNU_RECORDS")
	
	if [[ -z "$HTTP_STATUS" ]]; then
		bashio::log.error "Dynu Updater Initialisation - Invalid or empty JSON returned from Dynu Records API"
		exit 1
	fi
	
	EXCEPTION_TYPE=$(
	jq -r '.exception.type // empty' <<< "$DYNU_RECORDS")
	
	EXCEPTION_MESSAGE=$(
	jq -r '.exception.message // empty' <<< "$DYNU_RECORDS")

	case "$HTTP_STATUS" in
		200)
			bashio::log.info "Dynu Updater Initialisation - Dynu Records API data retrieved."
			;;
		*)
			bashio::log.error "Dynu Updater Initialisation - Dynu Records API Error: status: $HTTP_STATUS, type: $EXCEPTION_TYPE, message: $EXCEPTION_MESSAGE"
			exit 1
			;;

	esac
else
	bashio::log.error "Dynu Updater Initialisation - Dynu Records API could not be loaded"
	exit 1
fi

readarray -t DYNU_HOSTNAMES < <(
  jq -r '
	[.dnsRecords[]
	  | select(.recordType | IN("SOA","A","AAAA"))
	  | .hostname
	]
	| unique
	| sort_by( (split(".") | length), . )
	| .[]
  ' <<< "$DYNU_RECORDS"
)

if [[ ${#DYNU_HOSTNAMES[@]} -eq 0 ]]; then
	bashio::log.error "Dynu Updater Initialisation - no valid hostnames returned by Dynu Records API"
	exit 1
fi

IP_HOSTNAMES=()
CERT_HOSTNAMES=()

if [[ "$DNS_VALID" == "true" ]]; then
	for DOMAIN in $IP_DOMAINS; do
		if [[ "$DOMAIN" != "*.$DYNU_HOSTNAME" && ! " ${DYNU_HOSTNAMES[*]} " =~ " $DOMAIN " ]]; then
			bashio::log.error "Dynu Updater Initialisation - DNS Update Hostname $DOMAIN is not valid for this account"
			exit 1
		elif [[ "$DOMAIN" == "*.$DYNU_HOSTNAME" ]]; then
			# Warn if wildcard is not the only entry
			if [[ "$IP_DOMAINS" != "$DOMAIN" ]]; then
				bashio::log.warning "Dynu Updater Initialisation - Wildcard found, ignoring other DNS Update Hostnames"
			fi

			# Use all hostnames
			IP_HOSTNAMES=("${DYNU_HOSTNAMES[@]}")
			break
		fi
		IP_HOSTNAMES+=("$DOMAIN")
	done
	if [[ ${#IP_HOSTNAMES[@]} -eq 0 ]]; then
		bashio::log.error "Dynu Updater Initialisation - no valid hostnames for DNS Updater"
		exit 1
	fi
fi
if [[ "$CERT_VALID" == "true" ]]; then
	for DOMAIN in $CERT_DOMAINS; do
		if [[ "$DOMAIN" != "*.$DYNU_HOSTNAME" && ! " ${DYNU_HOSTNAMES[*]} " =~ " $DOMAIN " ]]; then
			bashio::log.error "Dynu Updater Initialisation - Certificate Hostname $DOMAIN is not valid for this account"
			exit 1
		elif [[ "$DOMAIN" == "*.$DYNU_HOSTNAME" ]]; then
			# Warn if wildcard is not the only entry
			if [[ "$CERT_DOMAINS" != "$DOMAIN" ]]; then
				bashio::log.warning "Dynu Updater Initialisation - Wildcard found, ignoring other Certificate Hostnames"
			fi

			# Use *.hostname only
			CERT_HOSTNAMES=("$DOMAIN")
			break
		fi
		CERT_HOSTNAMES+=("$DOMAIN")
	done
	if [[ ${#CERT_HOSTNAMES[@]} -eq 0 ]]; then
		bashio::log.error "Dynu Updater Initialisation - no valid hostnames for Certificate Manager"
		exit 1
	fi
fi

bashio::log.info "Dynu Updater Initialisation - Generating MQTT Configuration"
# MQTT connection

TLS_FLAGS=()
MQTT_ENTITIES=()

if [[ "$USE_MQTT" == "true" ]]; then
	if [[ -z "$MQTT_HOST" || ! "$MQTT_HOST" =~ ^([A-Za-z0-9._-]+|\[?[A-Fa-f0-9:]+\]?)$ \
		  || -z "$MQTT_PORT" || ! "$MQTT_PORT" =~ ^[0-9]+$ ]]; then
		bashio::log.error "Dynu Updater Initialisation - Invalid configuration for MQTT"
		if [[ -z "$MQTT_HOST" ]]; then
			bashio::log.error "  Hostname missing from config"
		elif [[ ! "$MQTT_HOST" =~ ^([A-Za-z0-9._-]+|\[?[A-Fa-f0-9:]+\]?)$ ]]; then
			bashio::log.error "  Hostname is not valid"
		fi
		if [[ -z "$MQTT_PORT" ]]; then
			bashio::log.error "  Port missing from config"
		elif [[ ! "$MQTT_PORT" =~ ^[0-9]+$ ]]; then
			bashio::log.error "  Port is not valid"
		fi
		exit 1
	fi
	if [[ "$MQTT_CORE" == "true" ]]; then
		if bashio::services.available 'mqtt'; then
			bashio::log.info "Dynu Updater Initialisation - Using Core Mosquitto Credentials"
			MQTT_USER=$(bashio::services mqtt username)
			MQTT_PASSWORD=$(bashio::services mqtt password)
		else
			bashio::log.warning "Dynu Updater Initialisation - Core Mosquitto Addon not installed"
		fi
	fi
	if [[ "$MQTT_CORE" == "false" ]] || ! bashio::services.available 'mqtt'; then
		if [[ -z "$MQTT_USER" || ! "$MQTT_USER" =~ ^[A-Za-z0-9._-]+$ \
			  || -z "$$MQTT_PASSWORD" || "$MQTT_PASSWORD" =~ [[:space:]] ]]; then
			bashio::log.error "Dynu Updater Initialisation - Invalid configuration for MQTT"
			if [[ -z "$MQTT_USER" ]]; then
				bashio::log.error "  User missing from config"
			elif [[ ! "$MQTT_USER" =~ ^[A-Za-z0-9._-]+$ ]]; then
				bashio::log.error "  User is not valid"
			fi
			if [[ -z "$MQTT_PASSWORD" ]]; then
				bashio::log.error "  Password missing from config"
			elif [[ "$MQTT_PASSWORD" =~ [[:space:]] ]]; then
				bashio::log.error "  Password is not valid"
			fi
			exit 1
		fi
	fi
	
	OPENSSL_FLAGS=()
	
	if [[ "$MQTT_USE_TLS" == "true" ]]; then
		if [[ -n "$MQTT_CAFILE" && ! -f "$MQTT_CAFILE" ]]; then
			bashio::log.error "Dynu Updater Initialisation - Specified CA File does not exist"
			exit 1
		fi
		if [[ -n "$MQTT_CERTFILE" && ! -f "$MQTT_CERTFILE" ]]; then
			bashio::log.error "Dynu Updater Initialisation - Specified Certificate File does not exist"
			exit 1
		fi
		if [[ -n "$MQTT_KEYFILE" && ! -f "$MQTT_KEYFILE" ]]; then
			bashio::log.error "Dynu Updater Initialisation - Specified Private Key File does not exist"
			exit 1
		fi
		if [[ -n "$MQTT_CAFILE" ]]; then
			TLS_FLAGS+=("--cafile" "$MQTT_CAFILE")
			OPENSSL_FLAGS+=("-CAfile" "$MQTT_CAFILE")
		else
			TLS_FLAGS+=("--tls-use-os-certs")
		fi
		if [[ -n "$MQTT_CERTFILE" ]]; then
			TLS_FLAGS+=("--cert" "$MQTT_CERTFILE")
			OPENSSL_FLAGS+=("-cert" "$MQTT_CERTFILE")
		fi
		if [[ -n "$MQTT_KEYFILE" ]]; then
			TLS_FLAGS+=("--key" "$MQTT_KEYFILE")
			OPENSSL_FLAGS+=("-key" "$MQTT_KEYFILE")
		fi
		if [[ "$MQTT_VALIDATE" == "false" ]]; then
			TLS_FLAGS+=("--insecure")
		else
			OPENSSL_FLAGS+=("-verify_hostname" "$MQTT_HOST")
		fi
	fi
	
	bashio::log.info "Dynu Updater Initialisation - Testing MQTT connection to $MQTT_HOST"
	bashio::log.info "Dynu Updater Initialisation - Testing MQTT connection with $MQTT_CAFILE"
	TEST_TOPIC="dynu_updater/selftest"
    TEST_PAYLOAD="selftest_$(date +%s)"
	if ! OUTPUT=$(mosquitto_pub \
        -h "$MQTT_HOST" \
        -u "$MQTT_USER" \
        -P "$MQTT_PASSWORD" \
        -p "$MQTT_PORT" \
        -t "$TEST_TOPIC" \
        -m "$TEST_PAYLOAD" \
        -q 1 \
		"${TLS_FLAGS[@]}" \
		2>&1); then

		bashio::log.error "Dynu Updater Initialisation - Unable to connect to $MQTT_HOST. $OUTPUT"
		
		if [[ "$MQTT_USE_TLS" == "true" ]]; then
			if ! TLS_OUTPUT=$(openssl s_client -connect "${MQTT_HOST}:${MQTT_PORT}" \
				  -servername "${MQTT_HOST}" \
				  "${OPENSSL_FLAGS[@]}" \
				  -verify 10 -verify_return_error < /dev/null 2>&1); then
				REASON=$(echo "$TLS_OUTPUT" | grep -E "verify error|Verify return code|Hostname mismatch|expired|unable|self signed" | head -n 1)
				bashio::log.error "TLS validation failed: $REASON"
			fi
		fi

		exit 1
	else
		bashio::log.info "Dynu Updater Initialisation - Connected to $MQTT_HOST"
	fi

	MQTT_ENTITIES=(
	  "last_ip_check|Last IP Check|timestamp|dynu_updater/last_ip_check|sensor||mdi:clock-check|dns"
	  "last_ip_update|Last IP Update|timestamp|dynu_updater/last_ip_update|sensor||mdi:update|dns"
	  "current_ipv4|Current IPv4||dynu_updater/current_ipv4|sensor||mdi:ip-network|dns"
	  "current_ipv6|Current IPv6||dynu_updater/current_ipv6|sensor||mdi:ip-network-outline|dns"
	  "next_cert_check|Next Certificate Check|timestamp|dynu_updater/next_cert_check|sensor||mdi:calendar-clock|cert"
	  "last_cert_check|Last Certificate Check|timestamp|dynu_updater/last_cert_check|sensor||mdi:certificate-outline|cert"
	  "last_cert_status|Certificate Valid|problem|dynu_updater/last_cert_status|binary_sensor|\"payload_on\": \"false\", \"payload_off\": \"true\",|mdi:shield-check|dns"
	  "wan_connection|WAN Connection|connectivity|dynu_updater/wan_connection|binary_sensor|\"payload_on\": \"true\", \"payload_off\": \"false\",|mdi:wan|dns"
	  "update_status|IP Update Status|problem|dynu_updater/update_status|binary_sensor|\"payload_on\": \"false\", \"payload_off\": \"true\",|mdi:sticker-check|dns"
	  "cert_expires|Certificate Expires|timestamp|dynu_updater/cert_expires|sensor||mdi:calendar-alert|cert"
	  "last_cert_renewal|Last Certificate Renewal|timestamp|dynu_updater/last_cert_renewal|sensor||mdi:certificate|cert"
	)

	# MQTT Device metadata
	DEVICE_ID="dynu_updater_$(hostname)"
	DEVICE_NAME="Dynu Updater"
	MANUFACTURER="TheUkShadow"
	MODEL="Dynu Updater Companion Integration"
	SW_VERSION=$(bashio::addon.version)

	# Shared availability topic
	AVAIL_TOPIC="dynu_updater/availability"

	publish_discovery() {
		if [[ "$USE_MQTT" == "true" ]]; then
			local entity_id="$1"
			local name="$2"
			local device_class="$3"
			local state_topic="$4"
			local component="$5"
			local extra="$6"
			local icon="$7"

			# Only include device_class if non-empty
			local device_class_json=""
			if [[ -n "$device_class" ]]; then
				device_class_json="\"device_class\": \"${device_class}\","
			fi
			
			local icon_json=""
			if [[ -n "$icon" ]]; then
				icon_json="\"icon\": \"${icon}\","
			fi
			
			# Per‑entity availability topic
			local entity_avail_topic="dynu_updater/${entity_id}/availability"

			if ! mosquitto_pub -h "$MQTT_HOST" -u "$MQTT_USER" -P "$MQTT_PASSWORD" -p "$MQTT_PORT" -r \
				  -t "homeassistant/${component}/${entity_id}/config" \
				  -m "{
					\"name\": \"${name}\",
					\"unique_id\": \"${DEVICE_ID}_${entity_id}\",
					\"state_topic\": \"${state_topic}\",
					\"availability\": [
					  { \"topic\": \"${AVAIL_TOPIC}\" },
					  { \"topic\": \"${entity_avail_topic}\" }
					],
					\"payload_available\": \"online\",
					\"payload_not_available\": \"offline\",
					${device_class_json}
					${icon_json}
					${extra}
					\"device\": {
					  \"identifiers\": [\"${DEVICE_ID}\"],
					  \"name\": \"${DEVICE_NAME}\",
					  \"manufacturer\": \"${MANUFACTURER}\",
					  \"model\": \"${MODEL}\",
					  \"sw_version\": \"${SW_VERSION}\"
					}
				  }" \
				  "${TLS_FLAGS[@]}"; then
				
				bashio::log.warning "Dynu Updater Initialisation - MQTT discovery failed for entity $entity_id"
			fi	
		fi
	}

	mqtt_publish() {
		if [[ "$USE_MQTT" == "true" ]]; then
			local topic="$1"
			local payload="$2"

			if ! mosquitto_pub \
				-h "$MQTT_HOST" \
				-p "$MQTT_PORT" \
				-u "$MQTT_USER" \
				-P "$MQTT_PASSWORD" \
				-t "dynu_updater/$topic" \
				-m "$payload" \
				-r \
				"${TLS_FLAGS[@]}"; then

				bashio::log.warning "Dynu Updater Initialisation - MQTT publish failed for topic $topic"
			fi
		fi
	}

	for entry in "${MQTT_ENTITIES[@]}"; do
		IFS="|" read -r entity_id name device_class state_topic component extra icon service <<< "$entry"

		publish_discovery \
			"$entity_id" \
			"$name" \
			"$device_class" \
			"$state_topic" \
			"$component" \
			"$extra" \
			"$icon" 

		# Publish initial per‑entity availability
		STATE="online"
		if [[ "$service" == "dns"  && "$DNS_VALID"  == "false" ]] ||
		   [[ "$service" == "cert" && "$CERT_VALID" == "false" ]]; then
			STATE="offline"
		fi
		
		mqtt_publish "${entity_id}/availability" "$STATE"
		
	done

	mqtt_publish "availability" "online"

fi
jq -n \
  --arg ipv6 "$DYNU_IPV6_ENABLED" \
  --argjson ip_domains "$(printf '%s\n' "${IP_HOSTNAMES[@]}" | jq -R . | jq -s .)" \
  --argjson cert_domains "$(printf '%s\n' "${CERT_HOSTNAMES[@]}" | jq -R . | jq -s .)" \
  --arg mqtt "$USE_MQTT" \
  --arg host "$MQTT_HOST" \
  --arg port "$MQTT_PORT" \
  --arg user "$MQTT_USER" \
  --arg password "$MQTT_PASSWORD" \
  --argjson tls_flags "$(printf '%s\n' "${TLS_FLAGS[@]}" | jq -R . | jq -s .)" \
  --argjson entities "$(printf '%s\n' "${MQTT_ENTITIES[@]}" | jq -R . | jq -s .)" \
  '{
      ipv6: $ipv6,
      ip_domains: $ip_domains,
	  cert_domains: $cert_domains,
	  mqtt: {
	    enabled: $mqtt,
	    host: $host,
	    port: $port,
	    user: $user,
	    password: $password,
		tls_flags: $tls_flags,
		entities: $entities
	  }
   }' \
  > /run/dynu-updater.conf


sleep 1

bashio::log.info "Dynu Updater Initialisation - Complete"
