# dynu-updater/rootfs/etc/s6-overlay/s6-rc.d/cert-manager/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start Certificate Manager
# ==============================================================================
set +e

KEYFILE=$(bashio::config 'keyfile' | xargs)
CERTFILE=$(bashio::config 'certfile' | xargs)
CERT_TIME=$(bashio::config 'update_cert_time' | xargs)
VERBOSE=$(bashio::config 'verbose')

CERT_DIR=/data/letsencrypt

if [[ -f /run/dynu-run-certbot.sh ]]; then

	while [ ! -f /run/dynu-dns-update.ready ]; do
		sleep 1
	done

	while true; do
		bashio::log.info "Certificate Manager - Running certificate renewal"
		
		# Run certbot
		if [[ "$VERBOSE" == "false" ]]; then
			/run/dynu-run-certbot.sh > /tmp/certbot.out 2>&1
			CERTBOT_EXIT=$?
		else
			/run/dynu-run-certbot.sh 2>&1 | tee /tmp/certbot.out
			CERTBOT_EXIT=${PIPESTATUS[0]}
		fi

		# Load output into a variable
		CERTBOT_OUTPUT=$(cat /tmp/certbot.out)
			
		# If certbot failed
		if [ "$CERTBOT_EXIT" -ne 0 ]; then
			bashio::log.error "Certificate Manager - Certbot failed with code $CERTBOT_EXIT"
			if [[ "$VERBOSE" == "false" ]]; then
				while IFS= read -r line; do
					bashio::log.error "  $line"
				done <<< "$CERTBOT_OUTPUT"
			fi
		else
			# If certbot succeeded, check if certificate was updated
			
			if echo "$CERTBOT_OUTPUT" | grep -q "Successfully received certificate" || true; then
				bashio::log.info "Certificate Manager - New certificate received"
				
				# Get the last modified cert directory and copy the cert and private key to store
				CERT_DIR_LATEST="$(ls -td $CERT_DIR/live/*/ | head -1)"
				cp "${CERT_DIR_LATEST}privkey.pem" "/ssl/$KEYFILE"
				cp "${CERT_DIR_LATEST}fullchain.pem" "/ssl/$CERTFILE"
				
				if bashio::config.has_value 'restart_addons'; then
					for index in $(bashio::config 'restart_addons|keys'); do
						sleep 1
						RAW_ADDON=$(bashio::config "restart_addons[${index}]")
						ADDON="${RAW_ADDON//-/_}"

						INFO_OUTPUT=$(curl --fail -s \
						  -H "Authorization: Bearer $HASSIO_TOKEN" \
						  http://hassio/addons/${ADDON}/info 2>&1)
						CURL_EXIT=$?

						if [ "$CURL_EXIT" -ne 0 ]; then
							bashio::log.warning "Certificate Manager - Failed to query add-on state: ${ADDON}"
							bashio::log.warning "  curl exit code: $CURL_EXIT"
							bashio::log.warning "  response: $INFO_OUTPUT"
							continue
						fi

						STATE=$(echo "$INFO_OUTPUT" | jq -r '.data.state')

						if [[ "$STATE" == "started" ]]; then

							bashio::log.info "Certificate Manager - Restarting add-on: ${ADDON}"

							RESTART_OUTPUT=$(curl --fail -s -X POST \
							  -H "Authorization: Bearer $HASSIO_TOKEN" \
							  -H "Content-Type: application/json" \
							  http://hassio/addons/${ADDON}/restart 2>&1)
							CURL_EXIT=$?

							if [ "$CURL_EXIT" -eq 0 ]; then
								bashio::log.info "Certificate Manager - Restarted add-on: ${ADDON}"
							else
								bashio::log.warning "Certificate Manager - Failed to restart add-on: ${ADDON}"
								bashio::log.warning "  curl exit code: $CURL_EXIT"
								bashio::log.warning "  response: $RESTART_OUTPUT"
							fi
						
						else
							bashio::log.warning "Certificate Manager - ${ADDON} is not running (state: $STATE)"
						fi
					done
				fi
			elif echo "$CERTBOT_OUTPUT" | grep -q "Certificate not yet due for renewal" || true; then
				bashio::log.info "Certificate Manager - Current certificate valid"
			elif [[ "$VERBOSE" == "false" ]]; then
				bashio::log.warning "Certificate Manager - Unexpected Certbot output:"
				while IFS= read -r line; do
					bashio::log.warning "  $line"
				done <<< "$CERTBOT_OUTPUT"
			fi
		fi

		now=$(date +%s)
		today_target=$(date -d "$(date +%Y-%m-%d) $CERT_TIME" +%s)
		
		# If the target time has already passed today, schedule for tomorrow
		if (( today_target <= now )); then
			next_run=$(date -d "tomorrow $CERT_TIME" +%s)
			TM="Tomorrow"
		else
			next_run=$today_target
			TM="Today"
		fi

		sleep_seconds=$(( next_run - now ))

		bashio::log.info "Certificate Manager - Sleeping until next run at $CERT_TIME ($TM)"
		sleep "$sleep_seconds"
	done
else
	sleep infinity
fi
