# dynu-updater/rootfs/etc/s6-overlay/s6-rc.d/server/run
#!/usr/bin/with-contenv bash

cd /usr/app

sleep 1

#exec python3 -u server.py