# dynu-updater/rootfs/etc/s6-overlay/s6-rc.d/server/run
#!/usr/bin/with-contenv bash
sleep 1
cd /usr/app
exec python3 -u server.py