# dynu-updater/rootfs/etc/s6-overlay/s6-rc.d/server/run
#!/usr/bin/with-contenv bashio
sleep 1
cd /usr/app
bashio::log.info "HTTP Server - Starting"
exec python3 -u server.py