# dynu-updater/rootfs/etc/s6-overlay/s6-rc.d/dns-update/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Dynamic DNS Updater
# ==============================================================================

HOSTNAME=$(bashio::config 'dynu_hostname' | tr '[:upper:]' '[:lower:]' | tr -d '[:space:]')
PASSWORD=$(bashio::config 'dynu_update_password' | xargs)
VERBOSE=$(bashio::config 'verbose')

update_status_json() {
    local key="$1"
    local value="$2"

    local final="/data/status.json"

    # Ensure file exists
    [ -f "$final" ] || echo "{}" > "$final"

    # Build updated JSON in memory
    local new_json
    new_json=$(jq --arg v "$value" --arg k "$key" '.[$k] = $v' "$final")

    # Compare old vs new
    if diff <(echo "$new_json") "$final" >/dev/null; then
        # No change → no write → no temp file
        return
    fi

    # Only now create a temp file
    local tmp
    tmp="$(mktemp /tmp/dynu_status.XXXXXX)"

    # Write new JSON to temp file
    echo "$new_json" > "$tmp"

    # Atomic replace
    mv "$tmp" "$final"
}

bashio::log.info "Dynu DDNS updater - starting"

if [[ -f /run/dynu-dns.conf ]]; then

	update_status_json "last_dns_check" "$(date -Is)"
	
	STATE_FILE="/run/dynu-dns.conf"

	DYNU_IPV6_ENABLED=$(jq -r '.ipv6' "$STATE_FILE")
	readarray -t DOMAIN_ARR < <(jq -r '.domains[]' "$STATE_FILE")
	
	STATUS_FILE="/data/status.json"

	if [ -f "$STATUS_FILE" ]; then
		LAST_IPV4=$(jq -r '.current_ipv4 // empty' "$STATUS_FILE")
		LAST_IPV6=$(jq -r '.current_ipv6 // empty' "$STATUS_FILE")
	else
		LAST_IPV4=""
		LAST_IPV6=""
	fi

	while true; do
		
		if [[ "$VERBOSE" == "true" ]]; then
			bashio::log.info "Dynu DDNS updater - Checking current IP addresses"
		fi
		
		# --- Detect IPv4 ---
		IPV4=$(curl -4 -s --fail --max-time 5 https://icanhazip.com | tr -d '[:space:]')

		if [[ -z "$IPV4" ]]; then
			bashio::log.error "Dynu DDNS updater - unable to confirm IPv4 address"
		fi
		
		if [[ -n "$IPV4" && ! "$IPV4" =~ ^([0-9]{1,3}\.){3}[0-9]{1,3}$ ]]; then
			bashio::log.error "Dynu DDNS updater - invalid IPv4 format returned: $IPV4"
			IPV4=""
		fi

		# --- Detect IPv6 only if domain supports it ---
		if [[ "$DYNU_IPV6_ENABLED" == "true" ]]; then
			IPV6=$(curl -6 -s --fail --max-time 5 https://icanhazip.com | tr -d '[:space:]')

			if [[ -z "$IPV6" ]]; then
				bashio::log.error "Dynu DDNS updater - unable to confirm IPv6 address"
			fi
			
			if [[ -n "$IPV6" && ! "$IPV6" =~ ^([0-9a-fA-F:]+)$ ]]; then
				bashio::log.error "Dynu DDNS updater - invalid IPv6 format returned: $IPV6"
				IPV6=""
			fi
			
		else
			IPV6=""
		fi

		if { [[ -n "$IPV4" && "$IPV4" != "$LAST_IPV4" ]]; } \
		   || { [[ -n "$IPV6" && "$IPV6" != "$LAST_IPV6" ]]; }; then

			bashio::log.info "Dynu DDNS updater - detected IP change"
			bashio::log.info "Dynu DDNS updater - New IPv4 = ${IPV4}"
			if [[ "$DYNU_IPV6_ENABLED" == "true" ]]; then
				bashio::log.info "Dynu DDNS updater - New IPv6 = ${IPV6}"
			fi

			for DOMAIN in "${DOMAIN_ARR[@]}"; do
				if [[ "$DOMAIN" == "$HOSTNAME" ]]; then
					bashio::log.info "Dynu DDNS updater - updating Host ${DOMAIN}"
					COMMAND="https://api.dynu.com/nic/update?hostname=${HOSTNAME}&password=${PASSWORD}"
				else
					bashio::log.info "Dynu DDNS updater - updating Alias ${DOMAIN}"
					ALIAS="${DOMAIN%.$HOSTNAME}"
					COMMAND="https://api.dynu.com/nic/update?hostname=${HOSTNAME}&alias=${ALIAS}&password=${PASSWORD}"
				fi

				if [[ -n "$IPV4" ]]; then
					COMMAND="${COMMAND}&myip=${IPV4}"
				fi

				if [[ -n "$IPV6" ]]; then
					COMMAND="${COMMAND}&myipv6=${IPV6}"
				fi

				RESPONSE="$(curl -s "${COMMAND}" || echo "No response")"

				bashio::log.info "Dynu DDNS updater - Server response = $RESPONSE"

			done
			
			if [[ ! -f /run/dynu-dns-update.ready ]]; then
				touch /run/dynu-dns-update.ready
			fi

			LAST_IPV4="$IPV4"
			LAST_IPV6="$IPV6"
			update_status_json "last_ip_update" "$(date -Is)"
			update_status_json "current_ipv4" "$IPV4"
			update_status_json "current_ipv6" "$IPV6"
			
		else
			if [[ "$VERBOSE" == "true" ]]; then
				bashio::log.info "Dynu DDNS updater - No update required"
			fi
		fi

		sleep 60
	done

else
	sleep infinity
fi
