# dynu_updater/rootfs/etc/s6-overlay/s6-rc.d/cert-manager/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Certificate Manager
# ==============================================================================

LOG_LEVEL=$(bashio::config 'log_level')
bashio::log.level "$LOG_LEVEL"

while [[ ! -f /run/dynu-dns-update.ready ]]; do
	sleep 1
done

if ! bashio::config.true 'cert.update_cert'; then
	bashio::log.info "Certificate Manager - Not required."
	sleep infinity
fi

DYNU_CONF="/run/dynu-updater.conf"

sleep 1

if [[ ! -f "$DYNU_CONF" ]]; then
	bashio::log.error "Certificate Manager - Dynu API Config not found."
	exit 1
fi

bashio::log.info "Certificate Manager - Started."

API_DATA=$(< "$DYNU_CONF")

DYNU_HOSTNAME=$(jq -r '.hostname' <<< "$API_DATA")

DYNU_API_KEY=$(jq -r '.cert.api_key' <<< "$API_DATA")
CHECK_DAYS=$(jq -r '.cert.cert_time' <<< "$API_DATA")
EMAIL=$(jq -r '.cert.email' <<< "$API_DATA")
KEYFILE=$(jq -r '.cert.keyfile' <<< "$API_DATA")
CERTFILE=$(jq -r '.cert.certfile' <<< "$API_DATA")
DRY_RUN=$(jq -r '.cert.dry_run' <<< "$API_DATA")
FORCE_RENEW=$(jq -r '.cert.force_renew' <<< "$API_DATA")

readarray -t DOMAIN_ARR < <(jq -r '.cert.domains[]' <<< "$API_DATA")

USE_MQTT=$(jq -r '.mqtt.enabled' <<< "$API_DATA")

if [[ "$USE_MQTT" == "true" ]]; then
	MQTT_HOST=$(jq -r '.mqtt.host' <<< "$API_DATA")
	MQTT_PORT=$(jq -r '.mqtt.port' <<< "$API_DATA")
	MQTT_USER=$(jq -r '.mqtt.user' <<< "$API_DATA")
	MQTT_PASSWORD=$(jq -r '.mqtt.password' <<< "$API_DATA")
	MQTT_CLIENT_ID=$(jq -r '.mqtt.client_id' <<< "$API_DATA")
	readarray -t TLS_FLAGS < <(jq -r '.mqtt.tls_flags[]' <<< "$API_DATA")
fi

CERT_DIR=/data/letsencrypt
CERT_PATH="${CERT_DIR}/live/$DYNU_HOSTNAME"
CERT_FILE="${CERT_PATH}/cert.pem"
WORK_DIR=/data/workdir
DNS_MULTI_CREDS=/run/dns-multi.ini
ACME_ARGUMENTS=()
KEY_ARGUMENTS=()
ADDITIONAL_ARGS=()
PROPAGATION_SECONDS=300

VALID_FROM_ISO=""
CREATED_ISO=""
EXPIRY_ISO=""

declare -A MQTT_RETRY_MAP=()

mqtt_publish() {
    if [[ "$USE_MQTT" == "true" ]]; then
        local topic="$1"
        local payload="$2"

        if ! mosquitto_pub \
			--quiet \
            -h "$MQTT_HOST" \
            -p "$MQTT_PORT" \
            -u "$MQTT_USER" \
            -P "$MQTT_PASSWORD" \
			-i "$MQTT_CLIENT_ID" \
            -t "dynu_updater/${topic}" \
            -m "$payload" \
            -r \
			"${TLS_FLAGS[@]}"; then

            bashio::log.warning "Certificate Manager - MQTT publish failed: Topic=${topic}."
			
			# Store only the latest payload for this topic
			MQTT_RETRY_MAP["$topic"]="$payload"
        fi
    fi
}

fire_event() {
	if bashio::config.true 'events'; then
		if ! curl -s --fail --max-time 5 -X POST \
		  -H "Authorization: Bearer $HASSIO_TOKEN" \
		  -H "Content-Type: application/json" \
		  -d "$1" \
		  "http://hassio/homeassistant/api/events/dynu_updater" >/dev/null; then
		  
			bashio::log.error "Dynu IP Updater - Could not send ${1} to Home Assistant Events."
		else
			bashio::log.debug "Dynu IP Updater - Sent ${1} to Home Assistant Events."
		fi
	fi
}

mkdir -p "${WORK_DIR}"
mkdir -p "${CERT_DIR}"

touch "${DNS_MULTI_CREDS}"
chmod 0600 "${DNS_MULTI_CREDS}"

# Start building the dns-multi credentials file
echo "# generated by letsencrypt app startup script" > "${DNS_MULTI_CREDS}"
echo "dns_multi_provider = dynu" >> "${DNS_MULTI_CREDS}"
echo "DYNU_API_KEY=${DYNU_API_KEY}" >> "${DNS_MULTI_CREDS}"
echo "DYNU_PROPAGATION_TIMEOUT=${PROPAGATION_SECONDS}" >> "${DNS_MULTI_CREDS}"

ACME_ARGUMENTS+=("--authenticator" "dns-multi" "--dns-multi-credentials" "${DNS_MULTI_CREDS}" "--dns-multi-propagation-seconds" "${PROPAGATION_SECONDS}")

ACME_ARGUMENTS+=("--preferred-chain" "ISRG Root X2")

# Set key-type
KEY_ARGUMENTS+=("--key-type" "ecdsa")
KEY_ARGUMENTS+=("--elliptic-curve" "secp256r1")

check_dates() {
	bashio::log.info "Certificate Manager - Checking Certificate dates."
	NOW_EPOCH=$(date +%s)
	VALID_FROM_RAW=""
	CREATED_RAW=""
	EXPIRY_RAW=""
	if [[ -f "$CERT_FILE" ]]; then
		if ! VALID_FROM_RAW=$(openssl x509 -in "${CERT_FILE}" -noout -startdate 2>/dev/null | cut -d= -f2); then
			VALID_FROM_RAW=""
		fi
		if ! CREATED_RAW=$(stat -c %y "${CERT_FILE}"); then
			CREATED_RAW=""
		fi
		if ! EXPIRY_RAW=$(openssl x509 -in "${CERT_FILE}" -noout -enddate 2>/dev/null | cut -d= -f2); then
			EXPIRY_RAW=""
		fi
		if [[ -n "$VALID_FROM_RAW" && -n "$CREATED_RAW" && -n "$EXPIRY_RAW" ]]; then
			VALID_FROM_ISO=$(date -u -d "${VALID_FROM_RAW}" +"%Y-%m-%dT%H:%M:%S+00:00")
			CREATED_ISO=$(date -u -d "${CREATED_RAW}" +"%Y-%m-%dT%H:%M:%S+00:00")
			EXPIRY_ISO=$(date -u -d "${EXPIRY_RAW}" +"%Y-%m-%dT%H:%M:%S+00:00")
			EXPIRY_EPOCH=$(date -d "${EXPIRY_RAW}" +%s)
			
			MAX_EPOCH=$(( NOW_EPOCH + (CHECK_DAYS * 86400) + 60 ))
			
			bashio::log.info "Certificate Manager - Certificate created ${CREATED_ISO}."
			bashio::log.info "Certificate Manager - Certificate valid from ${VALID_FROM_ISO}."
			bashio::log.info "Certificate Manager - Certificate expires ${EXPIRY_ISO}."
			mqtt_publish "cert_expires" "$EXPIRY_ISO"
            mqtt_publish "cert_created" "$CREATED_ISO"
			mqtt_publish "last_cert_check" "$(date -Is)"
			
			if (( NOW_EPOCH >= EXPIRY_EPOCH )); then
				bashio::log.warning "Certificate Manager - Certificate expired."
				mqtt_publish "last_cert_status" "false"
				NEXT_CHECK=$(( NOW_EPOCH + 600 ))
				return 1
			elif (( MAX_EPOCH >= EXPIRY_EPOCH )); then
				bashio::log.warning "Certificate Manager - Certificate will require renewing soon."
				mqtt_publish "last_cert_status" "true"
				NEXT_CHECK=$(( NOW_EPOCH + 600 ))
				return 0
			else
				bashio::log.info "Certificate Manager - Certificate valid."
				mqtt_publish "last_cert_status" "true"
				NEXT_CHECK=$(( EXPIRY_EPOCH - (CHECK_DAYS * 86400) ))
				if (( NOW_EPOCH >= NEXT_CHECK )); then
					NEXT_CHECK=$(( NOW_EPOCH + 600 ))
				fi
				return 0
			fi
		else
			bashio::log.error "Certificate Manager - Unable to parse certificate dates."
			mqtt_publish "last_cert_status" "false"
			NEXT_CHECK=$(( NOW_EPOCH + 600 ))
			return 1
		fi
	else
		bashio::log.error "Certificate Manager - Certificate file does not exist."
		mqtt_publish "last_cert_status" "false"
		NEXT_CHECK=$(( NOW_EPOCH + 600 ))
		return 1
	fi
}

check_cert() {
	bashio::log.info "Certificate Manager - Running Certbot. This can take up to ${PROPAGATION_SECONDS} seconds to complete if a new certificate is required."
	
	ADDITIONAL_ARGS=()
	
	if [[ "$FORCE_RENEW" == "true" ]]; then
		ADDITIONAL_ARGS+=("--force-renewal")
	fi

	if [[ "$DRY_RUN" == "true" ]]; then
		ADDITIONAL_ARGS+=("--dry-run")
	fi
	
	if [[ "$LOG_LEVEL" == "debug" ]]; then
		ADDITIONAL_ARGS+=("-vvv")
	fi
	
	# Start the dot-printer in the background
	(
		stdbuf -o0 -e0 bash -c '
			while true; do
				sleep 5
				printf "\033[32m.\033[0m"
			done
		'
	) &
	PROGRESS_PID=$!
	
	if CERTBOT_OUTPUT=$(certbot certonly --non-interactive --keep-until-expiring --expand \
		  --email "${EMAIL}" --agree-tos \
		  "${KEY_ARGUMENTS[@]}" \
		  "${ADDITIONAL_ARGS[@]}" \
		  --cert-name "${DYNU_HOSTNAME}" "${DOMAIN_ARR[@]/#/-d }" \
		  --config-dir "${CERT_DIR}" --work-dir "${WORK_DIR}" \
		  --preferred-challenges dns \
		  "${ACME_ARGUMENTS[@]}" \
		  2>&1); then
		CERTBOT_EXIT=0
	else
		CERTBOT_EXIT=$?
	fi
	
	# Stop the dot-printer
	kill "$PROGRESS_PID" 2>/dev/null
	echo ""

	bashio::log.debug "Certbot Output: ${CERTBOT_OUTPUT}"
	
	NOW_EPOCH=$(date +%s)
	
	# If certbot failed
	if [[ "$CERTBOT_EXIT" -ne 0 ]]; then
		bashio::log.error "Certificate Manager - Certbot failed with code ${CERTBOT_EXIT}."
		NEXT_CHECK=$(( NOW_EPOCH + 600 ))
	else
		# If certbot succeeded, check if certificate was updated
		if echo "$CERTBOT_OUTPUT" | grep -q "Successfully received certificate"; then
			bashio::log.info "Certificate Manager - New certificate received."
			
			# Copy the cert and private key to store
			cp "${CERT_PATH}/privkey.pem" "/ssl/$KEYFILE"
			cp "${CERT_PATH}/fullchain.pem" "/ssl/$CERTFILE"
			if check_dates; then
				fire_event "{\"action\": \"certificate_update\", \"status\": \"updated\", \"created\": \"${CREATED_ISO}\", \"expires\": \"${EXPIRY_ISO}\"}"
			else
				fire_event '{"action": "certificate_update", "status": "fail"}'
			fi
		elif echo "$CERTBOT_OUTPUT" | grep -q "Certificate not yet due for renewal"; then
			bashio::log.info "Certificate Manager - Current certificate did not require renewing."
			if check_dates; then
				fire_event "{\"action\": \"certificate_update\", \"status\": \"no_change\", \"created\": \"${CREATED_ISO}\", \"expires\": \"${EXPIRY_ISO}\"}"
			else
				fire_event '{"action": "certificate_update", "status": "fail"}'
			fi
		elif [[ "$DRY_RUN" == "true" ]] && echo "$CERTBOT_OUTPUT" | grep -q "The dry run was successful"; then	
			bashio::log.info "Certificate Manager - The dry-run was successful."
			if check_dates; then
				fire_event "{\"action\": \"certificate_update\", \"status\": \"no_change\", \"created\": \"${CREATED_ISO}\", \"expires\": \"${EXPIRY_ISO}\"}"
			else
				fire_event '{"action": "certificate_update", "status": "fail"}'
			fi
		else
			bashio::log.error "Certificate Manager - Unexpected Certbot output"
			NEXT_CHECK=$(( NOW_EPOCH + 600 ))
		fi
	fi
}

NEXT_CHECK=$(date +%s)

while true; do
	if [[ ${#MQTT_RETRY_MAP[@]} -ne 0 ]]; then
		for topic in "${!MQTT_RETRY_MAP[@]}"; do
            payload="${MQTT_RETRY_MAP[$topic]}"

            if mosquitto_pub \
				--quiet \
                -h "$MQTT_HOST" \
                -p "$MQTT_PORT" \
                -u "$MQTT_USER" \
                -P "$MQTT_PASSWORD" \
                -i "$MQTT_CLIENT_ID" \
                -t "dynu_updater/${topic}" \
                -m "$payload" \
                -r \
                "${TLS_FLAGS[@]}"; then

                bashio::log.info "Certificate Manager - Resent MQTT message: ${topic}"
                unset MQTT_RETRY_MAP["$topic"]
            fi
        done
	fi

	NOW_EPOCH=$(date +%s)
	
	if (( NOW_EPOCH >= NEXT_CHECK )); then
		check_cert
	
		NEXT_RUN_ISO=$(date -u -d "@$NEXT_CHECK" +"%Y-%m-%dT%H:%M:%S+00:00")
		bashio::log.info "Certificate Manager - Next renewal check: ${NEXT_RUN_ISO}."
		mqtt_publish "next_cert_check" "$NEXT_RUN_ISO"
		
		if [[ "FORCE_RENEW" == "true" || "DRY_RUN" == "true" ]]; then
			FORCE_RENEW="false"
			DRY_RUN="false"
		fi
	fi
	
	sleep 10
done

