# dynu_updater/rootfs/etc/s6-overlay/s6-rc.d/cert-manager/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start Certificate Manager
# ==============================================================================
set +e

if bashio::config.has_value 'dynu_hostname'; then
	DYNU_HOSTNAME=$(bashio::config 'dynu_hostname' | tr '[:upper:]' '[:lower:]' | tr -d '[:space:]')
else
	DYNU_HOSTNAME=""
fi
if bashio::config.has_value 'keyfile'; then
	KEYFILE=$(bashio::config 'keyfile' | xargs)
else
	KEYFILE="privkey.pem"
fi
if bashio::config.has_value 'certfile'; then
	CERTFILE=$(bashio::config 'certfile' | xargs)
else
	CERTFILE="fullchain.pem"
fi
if bashio::config.has_value 'update_cert_days'; then
	CHECK_DAYS=$(bashio::config 'update_cert_days')
else
	CHECK_DAYS=30
fi
if bashio::config.true 'verbose'; then
    VERBOSE="true"
else
    VERBOSE="false"
fi
CERT_DIR=/data/letsencrypt

LOOP_SLEEP=600

MQTT_CONF="/run/dynu-mqtt.conf"
USE_MQTT=$(jq -r '.mqtt' "$MQTT_CONF")
MQTT_HOST=$(jq -r '.host' "$MQTT_CONF")
MQTT_PORT=$(jq -r '.port' "$MQTT_CONF")
MQTT_USER=$(jq -r '.user' "$MQTT_CONF")
MQTT_PASSWORD=$(jq -r '.password' "$MQTT_CONF")

mqtt_publish() {
    if [[ "$USE_MQTT" == "true" ]]; then
        local topic="$1"
        local payload="$2"

        if ! mosquitto_pub \
            -h "$MQTT_HOST" \
            -p "$MQTT_PORT" \
            -u "$MQTT_USER" \
            -P "$MQTT_PASSWORD" \
            -t "dynu_updater/$topic" \
            -m "$payload" \
            -r; then

            bashio::log.warning "Certificate Manager - MQTT publish failed: Topic=$topic"
        fi
    fi
}

renew_cert() {
	if [[ "$VERBOSE" == "false" ]]; then
		/run/dynu-run-certbot.sh > /tmp/certbot.out 2>&1
		CERTBOT_EXIT=$?
	else
		/run/dynu-run-certbot.sh 2>&1 | tee /tmp/certbot.out
		CERTBOT_EXIT=${PIPESTATUS[0]}
	fi

	# Load output into a variable
	CERTBOT_OUTPUT=$(cat /tmp/certbot.out)
		
	# If certbot failed
	if [ "$CERTBOT_EXIT" -ne 0 ]; then
		bashio::log.error "Certificate Manager - Certbot failed with code $CERTBOT_EXIT"
		if [[ "$VERBOSE" == "false" ]]; then
			while IFS= read -r line; do
				bashio::log.error "  $line"
			done <<< "$CERTBOT_OUTPUT"
		fi
	else
		# If certbot succeeded, check if certificate was updated
		
		if echo "$CERTBOT_OUTPUT" | grep -q "Successfully received certificate"; then
			bashio::log.info "Certificate Manager - New certificate received"
			
			# Get the last modified cert directory and copy the cert and private key to store
			CERT_DIR_LATEST=$(ls -td "$CERT_DIR"/live/*/ 2>/dev/null | head -1)
			if [[ -z "$CERT_DIR_LATEST" ]]; then
				bashio::log.error "Certificate Manager - Unable to locate latest certificate directory"
				return
			fi
			cp "${CERT_DIR_LATEST}privkey.pem" "/ssl/$KEYFILE"
			cp "${CERT_DIR_LATEST}fullchain.pem" "/ssl/$CERTFILE"
			
			if bashio::config.has_value 'restart_addons'; then
				for index in $(bashio::config 'restart_addons|keys'); do
					sleep 1
					RAW_ADDON=$(bashio::config "restart_addons[${index}]")
					ADDON="${RAW_ADDON//-/_}"

					INFO_OUTPUT=$(curl --fail -s \
					  -H "Authorization: Bearer $HASSIO_TOKEN" \
					  http://hassio/addons/${ADDON}/info 2>&1)
					CURL_EXIT=$?

					if [ "$CURL_EXIT" -ne 0 ]; then
						bashio::log.warning "Certificate Manager - Failed to query add-on state: ${ADDON}"
						bashio::log.warning "  curl exit code: $CURL_EXIT"
						bashio::log.warning "  response: $INFO_OUTPUT"
						continue
					fi

					STATE=$(echo "$INFO_OUTPUT" | jq -r '.data.state')

					if [[ "$STATE" == "started" ]]; then

						bashio::log.info "Certificate Manager - Restarting add-on: ${ADDON}"

						RESTART_OUTPUT=$(curl --fail -s -X POST \
						  -H "Authorization: Bearer $HASSIO_TOKEN" \
						  -H "Content-Type: application/json" \
						  http://hassio/addons/${ADDON}/restart 2>&1)
						CURL_EXIT=$?

						if [ "$CURL_EXIT" -eq 0 ]; then
							bashio::log.info "Certificate Manager - Restarted add-on: ${ADDON}"
						else
							bashio::log.warning "Certificate Manager - Failed to restart add-on: ${ADDON}"
							bashio::log.warning "  curl exit code: $CURL_EXIT"
							bashio::log.warning "  response: $RESTART_OUTPUT"
						fi
					
					else
						bashio::log.warning "Certificate Manager - ${ADDON} is not running (state: $STATE)"
					fi
				done
			fi
		elif echo "$CERTBOT_OUTPUT" | grep -q "Certificate not yet due for renewal"; then
			bashio::log.info "Certificate Manager - Current certificate is valid"
		elif [[ "$VERBOSE" == "false" ]]; then
			bashio::log.warning "Certificate Manager - Unexpected Certbot output:"
			while IFS= read -r line; do
				bashio::log.warning "  $line"
			done <<< "$CERTBOT_OUTPUT"
		fi
	fi
}

check_cert_dates() {
    POST_CHECK="false"
    bashio::log.info "Certificate Manager - Checking certificate dates"
    mqtt_publish "last_cert_check" "$(date -Is)"

    CERT_PATH="/data/letsencrypt/live/$DYNU_HOSTNAME/cert.pem"

    if [[ -f "$CERT_PATH" ]]; then
        CREATED_RAW=$(openssl x509 -in "$CERT_PATH" -noout -startdate 2>/dev/null | cut -d= -f2)
        EXPIRY_RAW=$(openssl x509 -in "$CERT_PATH" -noout -enddate 2>/dev/null | cut -d= -f2)

        # If parsing fails, force renewal
        if [[ -z "$CREATED_RAW" || -z "$EXPIRY_RAW" ]]; then
            bashio::log.warning "Certificate Manager - Unable to parse certificate dates"
            mqtt_publish "last_cert_status" "false"
            renew_cert
            POST_CHECK="true"
        else
            CREATED_ISO=$(date -u -d "$CREATED_RAW" +"%Y-%m-%dT%H:%M:%S+00:00")
            EXPIRY_ISO=$(date -u -d "$EXPIRY_RAW" +"%Y-%m-%dT%H:%M:%S+00:00")

            mqtt_publish "cert_expires" "$EXPIRY_ISO"
            mqtt_publish "last_cert_renewal" "$CREATED_ISO"

            EXPIRY_EPOCH=$(date -d "$EXPIRY_RAW" +%s)
            NOW_EPOCH=$(date +%s)
            SEVEN_DAYS=$(( NOW_EPOCH + 7*24*60*60 ))
			
			bashio::log.info "Certificate Manager - Certificate expires $EXPIRY_ISO"
			
            if (( EXPIRY_EPOCH < NOW_EPOCH )); then
                bashio::log.info "Certificate Manager - Certificate expired, renewing"
				mqtt_publish "last_cert_status" "false"
                renew_cert
                POST_CHECK="true"
            elif (( EXPIRY_EPOCH <= SEVEN_DAYS )); then
				bashio::log.info "Certificate Manager - Certificate expires in less than 7 days, renewing"
                mqtt_publish "last_cert_status" "true"
                renew_cert
                POST_CHECK="true"
            else
				bashio::log.info "Certificate Manager - Certificate does not require renewing"
                RENEW_EPOCH=$(( EXPIRY_EPOCH - (CHECK_DAYS * 86400) ))
				LOOP_SLEEP=$(( RENEW_EPOCH - NOW_EPOCH ))
				if (( LOOP_SLEEP < 0 )); then
					LOOP_SLEEP=600
				fi
				mqtt_publish "last_cert_status" "true"
            fi
        fi
    else
        bashio::log.warning "Certificate Manager - Certificate missing/not yet created"
        mqtt_publish "last_cert_status" "false"
        renew_cert
        POST_CHECK="true"
    fi

    # Post-renewal verification
    if [[ "$POST_CHECK" == "true" ]]; then
        if [[ -f "$CERT_PATH" ]]; then
            CREATED_RAW=$(openssl x509 -in "$CERT_PATH" -noout -startdate 2>/dev/null | cut -d= -f2)
            EXPIRY_RAW=$(openssl x509 -in "$CERT_PATH" -noout -enddate 2>/dev/null | cut -d= -f2)

            if [[ -n "$CREATED_RAW" && -n "$EXPIRY_RAW" ]]; then
                CREATED_ISO=$(date -u -d "$CREATED_RAW" +"%Y-%m-%dT%H:%M:%S+00:00")
                EXPIRY_ISO=$(date -u -d "$EXPIRY_RAW" +"%Y-%m-%dT%H:%M:%S+00:00")

                mqtt_publish "cert_expires" "$EXPIRY_ISO"
                mqtt_publish "last_cert_renewal" "$CREATED_ISO"

                EXPIRY_EPOCH=$(date -d "$EXPIRY_RAW" +%s)
                NOW_EPOCH=$(date +%s)

                if (( EXPIRY_EPOCH < NOW_EPOCH )); then
                    mqtt_publish "last_cert_status" "false"
					LOOP_SLEEP=600
                else
					RENEW_EPOCH=$(( EXPIRY_EPOCH - (CHECK_DAYS * 86400) ))
					LOOP_SLEEP=$(( RENEW_EPOCH - NOW_EPOCH ))
					if (( LOOP_SLEEP < 0 )); then
						LOOP_SLEEP=600
					fi
                    mqtt_publish "last_cert_status" "true"
                fi
            else
                bashio::log.warning "Certificate Manager - Unable to parse renewed certificate dates."
                mqtt_publish "last_cert_status" "false"
				LOOP_SLEEP=600
            fi
        else
            bashio::log.warning "Certificate Manager - Certificate missing after renewal. Check the logs for more information"
            mqtt_publish "last_cert_status" "false"
			LOOP_SLEEP=600
        fi
    fi
}

if [[ -f /run/dynu-run-certbot.sh ]]; then

	while [ ! -f /run/dynu-dns-update.ready ]; do
		sleep 1
	done

	while true; do
		check_cert_dates
		
		NEXT_RUN_ISO=$(date -Is -d "@$(( $(date +%s) + LOOP_SLEEP ))")
		bashio::log.info "Certificate Manager - Next renewal check scheduled at $NEXT_RUN_ISO"
		mqtt_publish "next_cert_check" "$NEXT_RUN_ISO"
		
		sleep "$LOOP_SLEEP"
	done
else
	sleep infinity
fi
