# dynu_updater/rootfs/etc/s6-overlay/s6-rc.d/cert-manager/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Certificate Manager
# ==============================================================================

while [[ ! -f /run/dynu-dns-update.ready ]]; do
	sleep 1
done

if ! bashio::config.true 'update_cert'; then
	bashio::log.info "Certificate Manager - Not required"
	sleep infinity
fi

DYNU_HOSTNAME=$(bashio::config 'dynu_hostname' | tr '[:upper:]' '[:lower:]' | tr -d '[:space:]')
EMAIL=$(bashio::config 'email' | xargs)
KEYFILE=$(bashio::config 'keyfile' | xargs)
CERTFILE=$(bashio::config 'certfile' | xargs)
CHECK_DAYS=$(bashio::config 'update_cert_days')
DYNU_API_KEY=$(bashio::config 'dynu_api_key' | xargs)
CERT_DOMAINS=$(bashio::config 'domains' | tr '[:upper:]' '[:lower:]' | xargs)

if bashio::config.true 'verbose'; then
    VERBOSE="true"
else
    VERBOSE="false"
fi
if bashio::config.true 'force_renew'; then
    FORCE_RENEW="true"
else
    FORCE_RENEW="false"
fi
if bashio::config.true 'dry_run'; then
    DRY_RUN="true"
else
    DRY_RUN="false"
fi
if bashio::config.true 'mqtt'; then
    USE_MQTT="true"
else
    USE_MQTT="false"
fi

CERT_DIR=/data/letsencrypt
WORK_DIR=/data/workdir
# DNS_API_KEY=/data/dnsapikey
DNS_MULTI_CREDS=/data/dns-multi.ini
ACME_ARGUMENTS=()
KEY_ARGUMENTS=()
ADDITIONAL_ARGS=()
LOOP_SLEEP=600
POST_CHECK="false"
LAST_CHECK="ok"
PROPAGATION_SECONDS=300

MQTT_CONF="/run/dynu-mqtt.conf"

if [[ "$USE_MQTT" == "true" && -f "$MQTT_CONF" ]]; then
	MQTT_DATA=$(< "$MQTT_CONF")
	USE_MQTT=$(jq -r '.mqtt' <<< "$MQTT_DATA")
	MQTT_HOST=$(jq -r '.host' <<< "$MQTT_DATA")
	MQTT_PORT=$(jq -r '.port' <<< "$MQTT_DATA")
	MQTT_USER=$(jq -r '.user' <<< "$MQTT_DATA")
	MQTT_PASSWORD=$(jq -r '.password' <<< "$MQTT_DATA")
else
	USE_MQTT="false"
fi

mqtt_publish() {
    if [[ "$USE_MQTT" == "true" ]]; then
        local topic="$1"
        local payload="$2"

        if ! mosquitto_pub \
            -h "$MQTT_HOST" \
            -p "$MQTT_PORT" \
            -u "$MQTT_USER" \
            -P "$MQTT_PASSWORD" \
            -t "dynu_updater/$topic" \
            -m "$payload" \
            -r; then

            bashio::log.warning "Certificate Manager - MQTT publish failed: Topic=$topic"
        fi
    fi
}

# Call Dynu API
if ! DYNU_API_RESPONSE=$(curl -s --fail --max-time 5 -H "API-Key: $DYNU_API_KEY" "https://api.dynu.com/v2/dns" 2>/dev/null); then
	bashio::log.error "Certificate Manager - unable to contact Dynu API (network error or invalid API key)"
    sleep 1
	exit 1
fi

# Retrieve Domain ID
DYNU_DOMAIN_ID=$(
    jq -r --arg HOST "$DYNU_HOSTNAME" '
        .domains[]
        | select(.name == $HOST)
        | .id // empty
    ' <<< "$DYNU_API_RESPONSE")

if [[ -z "$DYNU_DOMAIN_ID" ]]; then
	bashio::log.error "Certificate Manager - '$DYNU_HOSTNAME' not found in Dynu domain list"
	sleep 1
	exit 1
fi

# Retrieve DNS Records for Domain
if ! DYNU_RECORDS=$(curl -s --fail --max-time 5 -H "API-Key: $DYNU_API_KEY" "https://api.dynu.com/v2/dns/${DYNU_DOMAIN_ID}/record" 2>/dev/null); then
	bashio::log.error "Certificate Manager - failed to retrieve valid DNS records from Dynu API"
	sleep 1
	exit 1
fi

DYNU_HOSTNAMES=$(
  jq -r '
    [.dnsRecords[]
      | select(.recordType | IN("SOA","A","AAAA"))
      | .hostname
    ]
    | unique
    | sort_by( (split(".") | length), . )
    | .[]
  ' <<< "$DYNU_RECORDS"
)

DOMAIN_ARR=()
for DOMAIN in $CERT_DOMAINS; do
	if [[ "$DOMAIN" == "*.$DYNU_HOSTNAME" ]]; then
		if [[ $(wc -w <<< "$CERT_DOMAINS") -ne 1 ]]; then
			bashio::log.warning "Certificate Manager - Wildcard found, ignoring other DNS Update Hostnames"
		fi
		DOMAIN_ARR=(-d "$DOMAIN")
		break
	elif ! echo "$DYNU_HOSTNAMES" | grep -qx "$DOMAIN"; then
		bashio::log.warning "Certificate Manager - $DOMAIN is not valid for this account"
		continue
	fi
	DOMAIN_ARR+=(-d "$DOMAIN")
done

if [[ ${#DOMAIN_ARR[@]} -eq 0 ]]; then
	bashio::log.error "Certificate Manager - no valid domains found"
	sleep 1
	exit 1
fi

mkdir -p "${WORK_DIR}"
mkdir -p "${CERT_DIR}"

# touch "${DNS_API_KEY}"
# chmod 0600 "${DNS_API_KEY}"
touch "${DNS_MULTI_CREDS}"
chmod 0600 "${DNS_MULTI_CREDS}"

# Start building the dns-multi credentials file
echo "# generated by letsencrypt app startup script" > "${DNS_MULTI_CREDS}"
echo "dns_multi_provider = dynu" >> "${DNS_MULTI_CREDS}"
echo "DYNU_API_KEY=$(bashio::config 'dynu_api_key')" >> "${DNS_MULTI_CREDS}"
echo "DYNU_PROPAGATION_TIMEOUT=${PROPAGATION_SECONDS}" >> "${DNS_MULTI_CREDS}"

ACME_ARGUMENTS+=("--authenticator" "dns-multi" "--dns-multi-credentials" "${DNS_MULTI_CREDS}" "--dns-multi-propagation-seconds" "${PROPAGATION_SECONDS}")

# Default Let's Encrypt servers.
ACME_ARGUMENTS+=("--preferred-chain" "ISRG Root X1")

if [[ "$FORCE_RENEW" == "true" ]]; then
	ADDITIONAL_ARGS+=("--force-renewal")
fi

if [[ "$VERBOSE" ==  "true" ]]; then
	ADDITIONAL_ARGS+=("-vvv")
fi

if [[ "$DRY_RUN" = "true" ]]; then
	ADDITIONAL_ARGS+=("--dry-run")
fi

# Set key-type
KEY_ARGUMENTS+=("--key-type" "ecdsa")

check_connectivity() {
    # Try Dynu (IPv4 or IPv6)
    if curl -4 -s --fail --max-time 5 https://checkip.dynu.com >/dev/null \
	  || curl -6 -s --fail --max-time 5 https://checkip.dynu.com >/dev/null; then
        return 0
    fi

    return 1
}

renew_cert() {
	POST_CHECK="true"
	until check_connectivity; do
        bashio::log.warning "Certificate Manager - Dynu not reachable, delaying certificate renewal"
        sleep 60
    done
	
	set -o pipefail
	
	TMP_OUTPUT=$(mktemp)

	if certbot certonly --non-interactive --keep-until-expiring --expand \
		  --email "$EMAIL" --agree-tos \
		  "${KEY_ARGUMENTS[@]}" \
		  "${ADDITIONAL_ARGS[@]}" \
		  --cert-name "$DYNU_HOSTNAME" "${DOMAIN_ARR[@]}" \
		  --config-dir "$CERT_DIR" --work-dir "$WORK_DIR" \
		  --preferred-challenges dns \
		  "${ACME_ARGUMENTS[@]}" \
		  2>&1 | tee "$TMP_OUTPUT"; then
		CERTBOT_EXIT=0
	else
		CERTBOT_EXIT=$?
	fi
	
	# If certbot failed
	if [[ "$CERTBOT_EXIT" -ne 0 ]]; then
		bashio::log.error "Certificate Manager - Certbot failed with code $CERTBOT_EXIT"
		LAST_CHECK="error"
	else
		# If certbot succeeded, check if certificate was updated
		# Load output into a variable
		CERTBOT_OUTPUT=$(< "$TMP_OUTPUT")
		
		if echo "$CERTBOT_OUTPUT" | grep -q "Successfully received certificate"; then
			bashio::log.info "Certificate Manager - New certificate received"
			
			# Get the last modified cert directory and copy the cert and private key to store
			CERT_DIR_LATEST=$(ls -td "$CERT_DIR"/live/*/ 2>/dev/null | head -1)
			if [[ -z "$CERT_DIR_LATEST" ]]; then
				bashio::log.error "Certificate Manager - Unable to locate latest certificate directory"
				return
			fi
			cp "${CERT_DIR_LATEST}privkey.pem" "/ssl/$KEYFILE"
			cp "${CERT_DIR_LATEST}fullchain.pem" "/ssl/$CERTFILE"
			
			if bashio::config.has_value 'restart_addons'; then
				for index in $(bashio::config 'restart_addons|keys'); do
					sleep 1
					RAW_ADDON=$(bashio::config "restart_addons[${index}]")
					ADDON="${RAW_ADDON//-/_}"

					if ! INFO_OUTPUT=$(curl --fail -s \
					  -H "Authorization: Bearer $HASSIO_TOKEN" \
					  http://hassio/addons/${ADDON}/info 2>&1); then
						bashio::log.warning "Certificate Manager - Failed to query state for addon: ${ADDON}"
						continue
					fi

					STATE=$(echo "$INFO_OUTPUT" | jq -r '.data.state')

					if [[ "$STATE" == "started" ]]; then

						bashio::log.info "Certificate Manager - Restarting add-on: ${ADDON}"

						if RESTART_OUTPUT=$(curl --fail -s -X POST \
						  -H "Authorization: Bearer $HASSIO_TOKEN" \
						  -H "Content-Type: application/json" \
						  http://hassio/addons/${ADDON}/restart 2>&1); then
							bashio::log.info "Certificate Manager - Restarted add-on: ${ADDON}"
						else
							bashio::log.warning "Certificate Manager - Failed to restart addon: ${ADDON}"
						fi
					
					else
						bashio::log.warning "Certificate Manager - ${ADDON} is not running (state: $STATE)"
					fi
				done
			fi
		elif echo "$CERTBOT_OUTPUT" | grep -q "Certificate not yet due for renewal"; then
			bashio::log.info "Certificate Manager - Current certificate is valid"
		elif [[ "$DRY_RUN" == "true" ]] && echo "$CERTBOT_OUTPUT" | grep -q "The dry run was successful"; then	
			bashio::log.info "Certificate Manager - The dry-run was successful"
		else
			bashio::log.warning "Certificate Manager - Unexpected Certbot output:"
		fi
	fi
}

check_cert_dates() {
    POST_CHECK="false"
	LAST_CHECK="ok"
    CERT_PATH="/data/letsencrypt/live/$DYNU_HOSTNAME/cert.pem"

	if [[ "$FORCE_RENEW" == "true" || "$DRY_RUN" == "true" ]]; then
		bashio::log.info "Certificate Manager - Skipping certificate date check: Force Renew/Dry Run"
        mqtt_publish "last_cert_status" "false"
        renew_cert
    elif [[ -f "$CERT_PATH" ]]; then
		bashio::log.info "Certificate Manager - Checking certificate dates"
		mqtt_publish "last_cert_check" "$(date -Is)"
        if ! CREATED_RAW=$(openssl x509 -in "$CERT_PATH" -noout -startdate 2>/dev/null | cut -d= -f2); then
			CREATED_RAW=""
		fi
		if ! EXPIRY_RAW=$(openssl x509 -in "$CERT_PATH" -noout -enddate 2>/dev/null | cut -d= -f2); then
			EXPIRY_RAW=""
		fi

        # If parsing fails, force renewal
        if [[ -z "$CREATED_RAW" || -z "$EXPIRY_RAW" ]]; then
            bashio::log.warning "Certificate Manager - Unable to parse certificate dates"
            mqtt_publish "last_cert_status" "false"
            renew_cert
        else
            CREATED_ISO=$(date -u -d "$CREATED_RAW" +"%Y-%m-%dT%H:%M:%S+00:00")
            EXPIRY_ISO=$(date -u -d "$EXPIRY_RAW" +"%Y-%m-%dT%H:%M:%S+00:00")

            mqtt_publish "cert_expires" "$EXPIRY_ISO"
            mqtt_publish "last_cert_renewal" "$CREATED_ISO"

            EXPIRY_EPOCH=$(date -d "$EXPIRY_RAW" +%s)
            NOW_EPOCH=$(date +%s)
            MAX_DAYS=$(( NOW_EPOCH + (CHECK_DAYS * 86400) + 60 ))
			
            if (( EXPIRY_EPOCH < NOW_EPOCH )); then
                bashio::log.info "Certificate Manager - Certificate expired, renewing"
				mqtt_publish "last_cert_status" "false"
                renew_cert
            elif (( EXPIRY_EPOCH <= MAX_DAYS )); then
				bashio::log.info "Certificate Manager - Certificate expires in less than $CHECK_DAYS days"
                mqtt_publish "last_cert_status" "true"
                renew_cert
            else
				bashio::log.info "Certificate Manager - Certificate expires $EXPIRY_ISO"
                RENEW_EPOCH=$(( EXPIRY_EPOCH - (CHECK_DAYS * 86400) ))
				LOOP_SLEEP=$(( RENEW_EPOCH - NOW_EPOCH ))
				if (( LOOP_SLEEP < 0 )); then
					LOOP_SLEEP=600
				fi
				mqtt_publish "last_cert_status" "true"
            fi
        fi
    else
        bashio::log.warning "Certificate Manager - Certificate missing/not yet created"
        mqtt_publish "last_cert_status" "false"
        renew_cert
    fi

    # Post-renewal verification
    if [[ "$POST_CHECK" == "true" ]]; then
        if [[ "$LAST_CHECK" == "error" ]]; then
			bashio::log.warning "Certificate Manager - Last renewal failed. Retrying in 10 minutes"
			LOOP_SLEEP=600
		elif [[ -f "$CERT_PATH" ]]; then
			if ! CREATED_RAW=$(openssl x509 -in "$CERT_PATH" -noout -startdate 2>/dev/null | cut -d= -f2); then
				CREATED_RAW=""
			fi
			if ! EXPIRY_RAW=$(openssl x509 -in "$CERT_PATH" -noout -enddate 2>/dev/null | cut -d= -f2); then
				EXPIRY_RAW=""
			fi

            if [[ -n "$CREATED_RAW" && -n "$EXPIRY_RAW" ]]; then
                CREATED_ISO=$(date -u -d "$CREATED_RAW" +"%Y-%m-%dT%H:%M:%S+00:00")
                EXPIRY_ISO=$(date -u -d "$EXPIRY_RAW" +"%Y-%m-%dT%H:%M:%S+00:00")

                mqtt_publish "cert_expires" "$EXPIRY_ISO"
                mqtt_publish "last_cert_renewal" "$CREATED_ISO"

                EXPIRY_EPOCH=$(date -d "$EXPIRY_RAW" +%s)
                NOW_EPOCH=$(date +%s)

                if (( EXPIRY_EPOCH < NOW_EPOCH )); then
                    mqtt_publish "last_cert_status" "false"
					LOOP_SLEEP=600
                else
					RENEW_EPOCH=$(( EXPIRY_EPOCH - (CHECK_DAYS * 86400) ))
					LOOP_SLEEP=$(( RENEW_EPOCH - NOW_EPOCH ))
					if (( LOOP_SLEEP < 0 )); then
						LOOP_SLEEP=600
					fi
                    mqtt_publish "last_cert_status" "true"
                fi
            else
                bashio::log.warning "Certificate Manager - Unable to parse renewed certificate dates."
                mqtt_publish "last_cert_status" "false"
				LOOP_SLEEP=600
            fi
        else
            bashio::log.warning "Certificate Manager - Certificate missing/not created after renewal. Check the logs for more information"
            mqtt_publish "last_cert_status" "false"
			LOOP_SLEEP=600
        fi
    fi
}

bashio::log.info "Certificate Manager - Started"

while true; do
	check_cert_dates
	
	NEXT_RUN_ISO=$(date -Is -d "@$(( $(date +%s) + LOOP_SLEEP ))")
	bashio::log.info "Certificate Manager - Next renewal check: $NEXT_RUN_ISO"
	mqtt_publish "next_cert_check" "$NEXT_RUN_ISO"
	
	sleep "$LOOP_SLEEP"
done

