# dynu_updater/rootfs/etc/s6-overlay/s6-rc.d/dns-update/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# IP Updater
# ==============================================================================

LOG_LEVEL=$(bashio::config 'log_level')
bashio::log.level "$LOG_LEVEL"

if ! bashio::config.true 'dns.update_ip'; then
	bashio::log.info "Dynu IP Updater - Not required."
	bashio::log.debug "Dynu IP Updater - Setting ready flag for Certificate Manager."
	touch /run/dynu-dns-update.ready
	sleep infinity
fi

DYNU_CONF="/run/dynu-updater.conf"

sleep 1

if [[ ! -f "$DYNU_CONF" ]]; then
	bashio::log.error "Dynu IP Updater - Dynu API Config not found."
	exit 1
fi

bashio::log.info "Dynu IP Updater - Started"

IPV6_REGEX='^((([0-9A-Fa-f]{1,4}:){7}[0-9A-Fa-f]{1,4})|(([0-9A-Fa-f]{1,4}:){1,7}:)|(([0-9A-Fa-f]{1,4}:){1,6}:[0-9A-Fa-f]{1,4})|(([0-9A-Fa-f]{1,4}:){1,5}(:[0-9A-Fa-f]{1,4}){1,2})|(([0-9A-Fa-f]{1,4}:){1,4}(:[0-9A-Fa-f]{1,4}){1,3})|(([0-9A-Fa-f]{1,4}:){1,3}(:[0-9A-Fa-f]{1,4}){1,4})|(([0-9A-Fa-f]{1,4}:){1,2}(:[0-9A-Fa-f]{1,4}){1,5})|([0-9A-Fa-f]{1,4}(:([0-9A-Fa-f]{1,4}:){1,6}))|(:((:[0-9A-Fa-f]{1,4}){1,7}|:))|(([0-9A-Fa-f]{1,4}:){1,4}((25[0-5]|2[0-4][0-9]|[01]?[0-9]?[0-9])\.){3}(25[0-5]|2[0-4][0-9]|[01]?[0-9]?[0-9])))$'
IPV4_REGEX='^([0-9]{1,3}\.){3}[0-9]{1,3}$'

API_DATA=$(< "$DYNU_CONF")

DYNU_HOSTNAME=$(jq -r '.hostname' <<< "$API_DATA")
PASSWORD=$(jq -r '.dns.password' <<< "$API_DATA")
IP_TIME=$(jq -r '.dns.ip_time' <<< "$API_DATA")
readarray -t IPV4_DOMAIN_ARR < <(jq -r '.dns.ipv4_domains[]' <<< "$API_DATA")
readarray -t IPV6_DOMAIN_ARR < <(jq -r '.dns.ipv6_domains[]' <<< "$API_DATA")

USE_MQTT=$(jq -r '.mqtt.enabled' <<< "$API_DATA")

if [[ "$USE_MQTT" == "true" ]]; then
	MQTT_HOST=$(jq -r '.mqtt.host' <<< "$API_DATA")
	MQTT_PORT=$(jq -r '.mqtt.port' <<< "$API_DATA")
	MQTT_USER=$(jq -r '.mqtt.user' <<< "$API_DATA")
	MQTT_PASSWORD=$(jq -r '.mqtt.password' <<< "$API_DATA")
	MQTT_CLIENT_ID=$(jq -r '.mqtt.client_id' <<< "$API_DATA")
	readarray -t TLS_FLAGS < <(jq -r '.mqtt.tls_flags[]' <<< "$API_DATA")
fi

declare -A MQTT_RETRY_MAP=()

mqtt_publish() {
    if [[ "$USE_MQTT" == "true" ]]; then
        local topic="$1"
        local payload="$2"

        if ! mosquitto_pub \
			--quiet \
            -h "$MQTT_HOST" \
            -p "$MQTT_PORT" \
            -u "$MQTT_USER" \
            -P "$MQTT_PASSWORD" \
			-i "$MQTT_CLIENT_ID" \
            -t "dynu_updater/${topic}" \
            -m "$payload" \
            -r \
			"${TLS_FLAGS[@]}"; then

            bashio::log.warning "Dynu IP Updater - MQTT publish failed: Topic=${topic}."
			
			# Store only the latest payload for this topic
			MQTT_RETRY_MAP["$topic"]="$payload"
        fi
    fi
}

fire_event() {
    if ! curl -s --fail --max-time 5 -X POST \
      -H "Authorization: Bearer $HASSIO_TOKEN" \
      -H "Content-Type: application/json" \
      -d "$1" \
      "http://hassio/homeassistant/api/events/dynu_updater" >/dev/null; then
	  
		bashio::log.error "Dynu IP Updater - Could not send ${1} to Home Assistant Events."
	else
		bashio::log.debug "Dynu IP Updater - Sent ${1} to Home Assistant Events."
	fi
}

IP_CONF="/data/dynu-ip.conf"

if [[ -f "$IP_CONF" ]]; then
	IP_DATA=$(< "$IP_CONF")
	CURRENT_IPV4=$(jq -r '.ipv4' <<< "$IP_DATA")
	CURRENT_IPV6=$(jq -r '.ipv6' <<< "$IP_DATA")
	LAST_IP_UPDATE=$(jq -r '.last' <<< "$IP_DATA")
else
	CURRENT_IPV4=""
	CURRENT_IPV6=""
	LAST_IP_UPDATE=""
fi

LAST_CHECK="true"
FIRST_RUN="true"

update_host() {
	bashio::log.info "Dynu IP Updater - Updating $1."
		
	if [[ "$1" == "$DYNU_HOSTNAME" ]]; then
		COMMAND="https://api.dynu.com/nic/update?hostname=${DYNU_HOSTNAME}&password=${PASSWORD}"
	else
		COMMAND="https://api.dynu.com/nic/update?hostname=${DYNU_HOSTNAME}&alias=${1%.$DYNU_HOSTNAME}&password=${PASSWORD}"
	fi
	
	if [[ -n "$2" ]]; then
		COMMAND="${COMMAND}&myip=${2}"
	fi

	if [[ -n "$3" ]]; then
		COMMAND="${COMMAND}&myipv6=${3}"
	fi
	
	bashio::log.debug "Dynu IP Updater - Update command = ${COMMAND}"
	
	if RESPONSE=$(curl -s --max-time 5 "${COMMAND}" 2>/dev/null); then
		case "$RESPONSE" in
			good*|nochg*)
				bashio::log.info "Dynu IP Updater - Success: ${RESPONSE}."
				return 0
				;;
			badauth)
				bashio::log.error "Dynu IP Updater - Authentication failed (badauth)."
				return 1
				;;
			notfqdn)
				bashio::log.error "Dynu IP Updater - Invalid hostname (notfqdn)."
				return 1
				;;
			nohost)
				bashio::log.error "Dynu IP Updater - Host not found (nohost)."
				return 1
				;;
			numhost)
				bashio::log.error "Dynu IP Updater - Too many hosts (numhost)."
				return 1
				;;
			abuse)
				bashio::log.error "Dynu IP Updater - Account blocked (abuse)."
				return 1
				;;
			*)
				bashio::log.error "Dynu IP Updater - Updating failed."
				return 1
				;;
		esac
	else
		bashio::log.error "Dynu IP Updater - Updating failed."
		return 1
	fi
}

check_ip_status() {
	bashio::log.debug "Dynu IP Updater - Checking current IP addresses."
	
	IPV4=""
	IPV6=""
	
	if [[ ${#IPV4_DOMAIN_ARR[@]} -ne 0 ]]; then
		# --- Detect IPv4 ---
		bashio::log.debug "Dynu IP Updater - Retrieving current IPv4 address from Dynu."
		if RAW=$(curl -4 -s --fail --max-time 5 https://checkip.dynu.com 2>/dev/null); then
			IPV4=$(echo "$RAW" | awk -F': ' '{print $2}' | tr -d '[:space:]')
			if [[ -z "$IPV4" || ! "$IPV4" =~ $IPV4_REGEX ]]; then
				bashio::log.warning "Dynu IP Updater - dynu.com returned an invalid IPv4 address."
				IPV4=""
			fi
		else
			bashio::log.warning "Dynu IP Updater - dynu.com IPv4 query failed."
		fi
		
		if [[ -z "$IPV4" || ! "$IPV4" =~ $IPV4_REGEX ]]; then
			bashio::log.debug "Dynu IP Updater - Retrieving current IPv4 address from Icanhazip."
			if RAW=$(curl -4 -s --fail --max-time 5 https://icanhazip.com 2>/dev/null); then
				IPV4=$(echo "$RAW" | tr -d '[:space:]')
				if [[ -z "$IPV4" || ! "$IPV4" =~ $IPV4_REGEX ]]; then
					bashio::log.warning "Dynu IP Updater - icanhazip.com returned an invalid IPv4 address."
					IPV4=""
				fi
			else
				bashio::log.warning "Dynu IP Updater - icanhazip.com IPv4 query failed."
			fi
		fi
		bashio::log.debug "Dynu IP Updater - Current IPv4 address: $IPV4."
	else
		bashio::log.debug "Dynu IP Updater - No IPv4 hostnames configured. Bypassing IPv4 check"
	fi
	
	# --- Detect IPv6 ---
	if [[ ${#IPV6_DOMAIN_ARR[@]} -ne 0 ]]; then
		bashio::log.debug "Dynu IP Updater - Retrieving current IPv6 address from Dynu."
		if RAW=$(curl -6 -s --fail --max-time 5 https://checkipv6.dynu.com 2>/dev/null); then
			IPV6=$(echo "$RAW" | awk -F': ' '{print $2}' | tr -d '[:space:]')
			if [[ -z "$IPV6" || ! "$IPV6" =~ $IPV6_REGEX ]]; then
				bashio::log.warning "Dynu IP Updater - dynu.com returned an invalid IPv6 address."
				IPV6=""
			fi
		else
			bashio::log.warning "Dynu IP Updater - dynu.com IPv6 query failed."
		fi
		
		if [[ -z "$IPV6" || ! "$IPV6" =~ $IPV6_REGEX ]]; then
			bashio::log.debug "Dynu IP Updater - Retrieving current IPv6 address from Icanhazip."
			if RAW=$(curl -6 -s --fail --max-time 5 https://icanhazip.com 2>/dev/null); then
				IPV6=$(echo "$RAW" | tr -d '[:space:]')
				if [[ -z "$IPV6" || ! "$IPV6" =~ $IPV6_REGEX ]]; then
					bashio::log.warning "Dynu IP Updater - icanhazip.com returned an invalid IPv6 address."
					IPV6=""
				fi
			else
				bashio::log.warning "Dynu IP Updater - icanhazip.com IPv6 query failed."
			fi
		fi
		bashio::log.debug "Dynu IP Updater - Current IPv6 address: $IPV6."
	else
		bashio::log.debug "Dynu IP Updater - No IPv6 hostnames configured. Bypassing IPv6 check"
	fi
	
	if [[ ${#IPV4_DOMAIN_ARR[@]} -eq 0 || -n "$IPV4" ]] && [[ ${#IPV6_DOMAIN_ARR[@]} -eq 0 || -n "$IPV6" ]]; then
		mqtt_publish "wan_connection" "true"
		mqtt_publish "last_ip_check" "$(date -Is)"
		fire_event "{\"action\": \"ip_check\", \"status\": \"ok\", \"data\": {\"ipv4\": \"${IPV4}\", \"ipv6\": \"${IPV6}\"}}"
		
		if [[ "$LAST_CHECK" == "false" || "$IPV4" != "$CURRENT_IPV4" || "$IPV6" != "$CURRENT_IPV6" ]]; then
			if [[ "$LAST_CHECK" == "false" ]]; then
				bashio::log.info "Dynu IP Updater - Last update failed. Updating all hosts."
			else
				bashio::log.info "Dynu IP Updater - IP address change detected. Updating all hosts."
			fi
			
			declare -A MERGED_ARRAY
			
			for HOST in "${IPV4_DOMAIN_ARR[@]}"; do
				MERGED_ARRAY["$HOST"]="$IPV4|"
			done
			
			for HOST in "${IPV6_DOMAIN_ARR[@]}"; do
				if [[ -n "${MERGED_ARRAY[$HOST]}" ]]; then
					MERGED_ARRAY["$HOST"]="${MERGED_ARRAY[$HOST]}$IPV6"
				else
					MERGED_ARRAY["$HOST"]="|$IPV6"
				fi
			done
			
			UPDATED=${#MERGED_ARRAY[@]}
			
			for HOST in "${!MERGED_ARRAY[@]}"; do
				IFS='|' read -r IP4 IP6 <<< "${MERGED_ARRAY[$HOST]}"
				if update_host "$HOST" "$IP4" "$IP6"; then
					UPDATED=$((UPDATED - 1))
				fi
			done
	
			if (( UPDATED > 0 )); then
				bashio::log.error "Dynu IP Updater - ${UPDATED} Hosts/Aliases could not be updated."
				LAST_CHECK="false"
				mqtt_publish "update_status" "false"
				NOW_EPOCH=$(date +%s)
				NEXT_CHECK=$(( NOW_EPOCH + 60 ))
				fire_event '{"action": "ip_update", "status": "fail"}'
			else
				CURRENT_IPV4="$IPV4"
				CURRENT_IPV6="$IPV6"
				LAST_IP_UPDATE="$(date -Is)"
				LAST_CHECK="true"
				
				mqtt_publish "update_status" "true"
				mqtt_publish "current_ipv4" "$IPV4"
				mqtt_publish "current_ipv6" "$IPV6"
				mqtt_publish "last_ip_update" "$(date -Is)"
				
				jq -n \
				  --arg js_ipv4 "$CURRENT_IPV4" \
				  --arg js_ipv6 "$CURRENT_IPV6" \
				  --arg js_last "$(date -Is)" \
				  '{
					  ipv4: $js_ipv4,
					  ipv6: $js_ipv6,
					  last: $js_last
				   }' \
				  > "$IP_CONF"
				  
				bashio::log.info "Dynu IP Updater - All Hosts/Aliases were updated."
				NOW_EPOCH=$(date +%s)
				NEXT_CHECK=$(( NOW_EPOCH + (IP_TIME * 60) ))
				fire_event '{"action": "ip_update", "status": "updated"}'
			fi
		else
			bashio::log.debug "Dynu IP Updater - No IP address change detected."
			NOW_EPOCH=$(date +%s)
			NEXT_CHECK=$(( NOW_EPOCH + (IP_TIME * 60) ))
			fire_event '{"action": "ip_update", "status": "no_change"}'
		fi

	else
		mqtt_publish "wan_connection" "false"
		bashio::log.error "Dynu IP Updater - Unable to retrieve current IP address data."
		NOW_EPOCH=$(date +%s)
		NEXT_CHECK=$(( NOW_EPOCH + 60 ))
		fire_event '{"action": "ip_check", "status": "fail"}'
	fi
	
	if [[ "$FIRST_RUN" == "true" ]]; then
		bashio::log.debug "Dynu IP Updater - Setting ready flag for Certificate Manager."
		touch /run/dynu-dns-update.ready
		FIRST_RUN="false"
	fi
	
	bashio::log.debug "Dynu IP Updater - Checking and updating complete."
}

NEXT_CHECK=$(date +%s)

while true; do
	if [[ ${#MQTT_RETRY_MAP[@]} -ne 0 ]]; then
		for topic in "${!MQTT_RETRY_MAP[@]}"; do
            payload="${MQTT_RETRY_MAP[$topic]}"
			
            if mosquitto_pub \
				--quiet \
                -h "$MQTT_HOST" \
                -p "$MQTT_PORT" \
                -u "$MQTT_USER" \
                -P "$MQTT_PASSWORD" \
                -i "$MQTT_CLIENT_ID" \
                -t "dynu_updater/${topic}" \
                -m "$payload" \
                -r \
                "${TLS_FLAGS[@]}"; then

                bashio::log.info "Dynu IP Updater - Resent MQTT message: ${topic}"
                unset MQTT_RETRY_MAP["$topic"]
            fi
        done
	fi
	
	NOW_EPOCH=$(date +%s)
	
	if (( NOW_EPOCH >= NEXT_CHECK )); then
		check_ip_status
	fi

	sleep 10
done