# dynu_updater/rootfs/etc/s6-overlay/s6-rc.d/dns-update/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# DNS Updater
# ==============================================================================

if ! bashio::config.true 'dns.update_ip'; then
	bashio::log.info "Dynu DNS Updater - Not required."
	touch /run/dynu-dns-update.ready
	sleep infinity
fi

DYNU_CONF="/run/dynu-updater.conf"

sleep 1

if [[ ! -f "$DYNU_CONF" ]]; then
	bashio::log.error "Dynu DNS Updater - Dynu API Config not found."
	exit 1
fi

bashio::log.info "Dynu DNS Updater - Started"

IPV6_REGEX='^((([0-9A-Fa-f]{1,4}:){7}[0-9A-Fa-f]{1,4})|(([0-9A-Fa-f]{1,4}:){1,7}:)|(([0-9A-Fa-f]{1,4}:){1,6}:[0-9A-Fa-f]{1,4})|(([0-9A-Fa-f]{1,4}:){1,5}(:[0-9A-Fa-f]{1,4}){1,2})|(([0-9A-Fa-f]{1,4}:){1,4}(:[0-9A-Fa-f]{1,4}){1,3})|(([0-9A-Fa-f]{1,4}:){1,3}(:[0-9A-Fa-f]{1,4}){1,4})|(([0-9A-Fa-f]{1,4}:){1,2}(:[0-9A-Fa-f]{1,4}){1,5})|([0-9A-Fa-f]{1,4}(:([0-9A-Fa-f]{1,4}:){1,6}))|(:((:[0-9A-Fa-f]{1,4}){1,7}|:))|(([0-9A-Fa-f]{1,4}:){1,4}((25[0-5]|2[0-4][0-9]|[01]?[0-9]?[0-9])\.){3}(25[0-5]|2[0-4][0-9]|[01]?[0-9]?[0-9])))$'
IPV4_REGEX='^([0-9]{1,3}\.){3}[0-9]{1,3}$'

API_DATA=$(< "$DYNU_CONF")

DYNU_HOSTNAME=$(jq -r '.hostname' <<< "$API_DATA")
VERBOSE=$(jq -r '.verbose' <<< "$API_DATA")
PASSWORD=$(jq -r '.dns.password' <<< "$API_DATA")
DYNU_IPV6_ENABLED=$(jq -r '.dns.ipv6' <<< "$API_DATA")
readarray -t DOMAIN_ARR < <(jq -r '.dns.domains[]' <<< "$API_DATA")

USE_MQTT=$(jq -r '.mqtt.enabled' <<< "$API_DATA")

if [[ "$USE_MQTT" == "true" ]]; then
	MQTT_HOST=$(jq -r '.mqtt.host' <<< "$API_DATA")
	MQTT_PORT=$(jq -r '.mqtt.port' <<< "$API_DATA")
	MQTT_USER=$(jq -r '.mqtt.user' <<< "$API_DATA")
	MQTT_PASSWORD=$(jq -r '.mqtt.password' <<< "$API_DATA")
	MQTT_CLIENT_ID=$(jq -r '.mqtt.client_id' <<< "$API_DATA")
	readarray -t TLS_FLAGS < <(jq -r '.mqtt.tls_flags[]' <<< "$API_DATA")
fi

declare -A MQTT_RETRY_MAP=()

mqtt_publish() {
    if [[ "$USE_MQTT" == "true" ]]; then
        local topic="$1"
        local payload="$2"

        if ! mosquitto_pub \
			--quiet \
            -h "$MQTT_HOST" \
            -p "$MQTT_PORT" \
            -u "$MQTT_USER" \
            -P "$MQTT_PASSWORD" \
			-i "$MQTT_CLIENT_ID" \
            -t "dynu_updater/${topic}" \
            -m "$payload" \
            -r \
			"${TLS_FLAGS[@]}"; then

            bashio::log.warning "Dynu DNS Updater - MQTT publish failed: Topic=${topic}."
			
			# Store only the latest payload for this topic
			MQTT_RETRY_MAP["$topic"]="$payload"
        fi
    fi
}

IP_CONF="/data/dynu-ip.conf"

if [[ -f "$IP_CONF" ]]; then
	IP_DATA=$(< "$IP_CONF")
	CURRENT_IPV4=$(jq -r '.ipv4' <<< "$IP_DATA")
	CURRENT_IPV6=$(jq -r '.ipv6' <<< "$IP_DATA")
	LAST_IP_UPDATE=$(jq -r '.last' <<< "$IP_DATA")
else
	CURRENT_IPV4=""
	CURRENT_IPV6=""
	LAST_IP_UPDATE=""
fi

LAST_CHECK="true"
FIRST_RUN="true"

update_host() {
	bashio::log.info "Dynu DNS Updater - Updating $1."
				
	if [[ "$1" == "$DYNU_HOSTNAME" ]]; then
		COMMAND="https://api.dynu.com/nic/update?hostname=${DYNU_HOSTNAME}&password=${PASSWORD}&myip=${2}"
	else
		COMMAND="https://api.dynu.com/nic/update?hostname=${DYNU_HOSTNAME}&alias=${1%.$DYNU_HOSTNAME}&password=${PASSWORD}&myip=${2}"
	fi

	if [[ -n "$3" ]]; then
		COMMAND="${COMMAND}&myipv6=${3}"
	fi
	
	if RESPONSE=$(curl -s --max-time 5 "${COMMAND}" 2>/dev/null); then
		case "$RESPONSE" in
			good*|nochg*)
				bashio::log.info "Dynu DNS Updater - Success: ${RESPONSE}."
				return 0
				;;
			badauth)
				bashio::log.error "Dynu DNS Updater - Authentication failed (badauth)."
				return 1
				;;
			notfqdn)
				bashio::log.error "Dynu DNS Updater - Invalid hostname (notfqdn)."
				return 1
				;;
			nohost)
				bashio::log.error "Dynu DNS Updater - Host not found (nohost)."
				return 1
				;;
			numhost)
				bashio::log.error "Dynu DNS Updater - Too many hosts (numhost)."
				return 1
				;;
			abuse)
				bashio::log.error "Dynu DNS Updater - Account blocked (abuse)."
				return 1
				;;
			*)
				bashio::log.error "Dynu DNS Updater - Updating failed."
				return 1
				;;
		esac
	else
		bashio::log.error "Dynu DNS Updater - Updating failed."
		return 1
	fi
}

mkfifo /tmp/mqtt_pipe
mosquitto_sub \
	-h "$MQTT_HOST" \
	-p "$MQTT_PORT" \
	-u "$MQTT_USER" \
	-P "$MQTT_PASSWORD" \
	-i "${MQTT_CLIENT_ID}_listener" \
	-t "dynu_updater/ip/service" \
	"${TLS_FLAGS[@]}" > /tmp/mqtt_pipe &

bashio::log.info "Starting loop"

NEXT_CHECK=$(date +%s)

while true; do
	if [[ ${#MQTT_RETRY_MAP[@]} -ne 0 ]]; then
		for topic in "${!MQTT_RETRY_MAP[@]}"; do
            payload="${MQTT_RETRY_MAP[$topic]}"
			
            if mosquitto_pub \
				--quiet \
                -h "$MQTT_HOST" \
                -p "$MQTT_PORT" \
                -u "$MQTT_USER" \
                -P "$MQTT_PASSWORD" \
                -i "$MQTT_CLIENT_ID" \
                -t "dynu_updater/${topic}" \
                -m "$payload" \
                -r \
                "${TLS_FLAGS[@]}"; then

                bashio::log.info "Dynu DNS Updater - Resent MQTT message: ${topic}"
                unset MQTT_RETRY_MAP["$topic"]
            fi
        done
	fi
	
	if read -t 0.1 -r line < /tmp/mqtt_pipe; then
        bashio::log.info "$line"
    fi
	
	NOW_EPOCH=$(date +%s)
	
	if (( NOW_EPOCH >= NEXT_CHECK )); then
	
		if [[ "$VERBOSE" == "true" ]]; then
			bashio::log.info "Dynu DNS Updater - Checking current IP addresses."
		fi
		
		IPV4=""
		IPV6=""
		
		# --- Detect IPv4 ---
		if [[ "$VERBOSE" == "true" ]]; then
			bashio::log.info "Dynu DNS Updater - Retrieving current IPv4 address from Dynu."
		fi
		if RAW=$(curl -4 -s --fail --max-time 5 https://checkip.dynu.com 2>/dev/null); then
			IPV4=$(echo "$RAW" | awk -F': ' '{print $2}' | tr -d '[:space:]')
			if [[ -z "$IPV4" || ! "$IPV4" =~ $IPV4_REGEX ]]; then
				bashio::log.warning "Dynu DNS Updater - dynu.com returned an invalid IPv4 address."
				IPV4=""
			fi
		else
			bashio::log.warning "Dynu DNS Updater - dynu.com IPv4 query failed."
		fi
		
		if [[ -z "$IPV4" || ! "$IPV4" =~ $IPV4_REGEX ]]; then
			if [[ "$VERBOSE" == "true" ]]; then
				bashio::log.info "Dynu DNS Updater - Retrieving current IPv4 address from Icanhazip."
			fi
			if RAW=$(curl -4 -s --fail --max-time 5 https://icanhazip.com 2>/dev/null); then
				IPV4=$(echo "$RAW" | tr -d '[:space:]')
				if [[ -z "$IPV4" || ! "$IPV4" =~ $IPV4_REGEX ]]; then
					bashio::log.warning "Dynu DNS Updater - icanhazip.com returned an invalid IPv4 address."
					IPV4=""
				fi
			else
				bashio::log.warning "Dynu DNS Updater - icanhazip.com IPv4 query failed."
			fi
		fi
		
		# --- Detect IPv6 only if domain supports it ---
		if [[ "$DYNU_IPV6_ENABLED" == "true" ]]; then
			if [[ "$VERBOSE" == "true" ]]; then
				bashio::log.info "Dynu DNS Updater - Retrieving current IPv6 address from Dynu."
			fi
			if RAW=$(curl -6 -s --fail --max-time 5 https://checkipv6.dynu.com 2>/dev/null); then
				IPV6=$(echo "$RAW" | awk -F': ' '{print $2}' | tr -d '[:space:]')
				if [[ -z "$IPV6" || ! "$IPV6" =~ $IPV6_REGEX ]]; then
					bashio::log.warning "Dynu DNS Updater - dynu.com returned an invalid IPv6 address."
					IPV6=""
				fi
			else
				bashio::log.warning "Dynu DNS Updater - dynu.com IPv6 query failed."
			fi
			
			if [[ -z "$IPV6" || ! "$IPV6" =~ $IPV6_REGEX ]]; then
				if [[ "$VERBOSE" == "true" ]]; then
					bashio::log.info "Dynu DNS Updater - Retrieving current IPv6 address from Icanhazip."
				fi
				if RAW=$(curl -6 -s --fail --max-time 5 https://icanhazip.com 2>/dev/null); then
					IPV6=$(echo "$RAW" | tr -d '[:space:]')
					if [[ -z "$IPV6" || ! "$IPV6" =~ $IPV6_REGEX ]]; then
						bashio::log.warning "Dynu DNS Updater - icanhazip.com returned an invalid IPv6 address."
						IPV6=""
					fi
				else
					bashio::log.warning "Dynu DNS Updater - icanhazip.com IPv6 query failed."
				fi
			fi
		fi
		
		if [[ "$VERBOSE" == "true" ]]; then
			bashio::log.info "Dynu DNS Updater - Current IPv4 address: $IPV4."
			if [[ "$DYNU_IPV6_ENABLED" == "true" ]]; then
				bashio::log.info "Dynu DNS Updater - Current IPv6 address: $IPV6."
			fi
		fi
		
		if [[ -n "$IPV4" ]] && [[ "$DYNU_IPV6_ENABLED" == "false" || -n "$IPV6" ]]; then
			mqtt_publish "wan_connection" "true"
			mqtt_publish "last_ip_check" "$(date -Is)"
			
			if [[ "$LAST_CHECK" == "false" || "$IPV4" != "$CURRENT_IPV4" || "$IPV6" != "$CURRENT_IPV6" ]]; then
				if [[ "$LAST_CHECK" == "false" ]]; then
					bashio::log.info "Dynu DNS Updater - Last update failed. Updating all hosts."
				else
					bashio::log.info "Dynu DNS Updater - IP address change detected. Updating all hosts."
				fi
							  
				UPDATED=${#DOMAIN_ARR[@]}
				
				for HOST in "${DOMAIN_ARR[@]}"; do
					if update_host "$HOST" "$IPV4" "$IPV6"; then
						UPDATED=$((UPDATED - 1))
					fi
				done
				
				if (( UPDATED > 0 )); then
					bashio::log.error "Dynu DNS Updater - ${UPDATED} Hosts/Aliases could not be updated."
					LAST_CHECK="false"
					mqtt_publish "update_status" "false"
				else
					CURRENT_IPV4="$IPV4"
					CURRENT_IPV6="$IPV6"
					LAST_IP_UPDATE="$(date -Is)"
					LAST_CHECK="true"
					
					mqtt_publish "update_status" "true"
					mqtt_publish "current_ipv4" "$IPV4"
					mqtt_publish "current_ipv6" "$IPV6"
					mqtt_publish "last_ip_update" "$(date -Is)"
					
					jq -n \
					  --arg js_ipv4 "$CURRENT_IPV4" \
					  --arg js_ipv6 "$CURRENT_IPV6" \
					  --arg js_last "$(date -Is)" \
					  '{
						  ipv4: $js_ipv4,
						  ipv6: $js_ipv6,
						  last: $js_last
					   }' \
					  > "$IP_CONF"
					  
					bashio::log.info "Dynu DNS Updater - All Hosts/Aliases were updated."
				fi
			elif [[ "$VERBOSE" == "true" || "$FIRST_RUN" == "true" ]]; then
				bashio::log.info "Dynu DNS Updater - No IP address change detected."
			fi

		else
			mqtt_publish "wan_connection" "false"
			bashio::log.error "Dynu DNS Updater - Unable to retrieve current IP address data."
		fi
		
		if [[ "$FIRST_RUN" == "true" ]]; then
			if [[ "$VERBOSE" == "true" ]]; then
				bashio::log.info "Dynu DNS Updater - Setting ready flag for Certificate Manager."
			fi
			touch /run/dynu-dns-update.ready
			FIRST_RUN="false"
		fi
		
		if [[ "$VERBOSE" == "true" ]]; then
			bashio::log.info "Dynu DNS Updater - Checking and updating complete."
		fi
		
		NOW_EPOCH=$(date +%s)
		NEXT_CHECK=$(( NOW_EPOCH + 60 ))
	fi

	sleep 1
done