# dynu_updater/rootfs/etc/s6-overlay/s6-rc.d/dns-update/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# DNS Updater
# ==============================================================================

if ! bashio::config.true 'update_ip'; then
	bashio::log.info "Dynu DNS Updater - Not required"
	touch /run/dynu-dns-update.ready
	sleep infinity
fi

IPV6_REGEX='^((([0-9A-Fa-f]{1,4}:){7}[0-9A-Fa-f]{1,4})|(([0-9A-Fa-f]{1,4}:){1,7}:)|(([0-9A-Fa-f]{1,4}:){1,6}:[0-9A-Fa-f]{1,4})|(([0-9A-Fa-f]{1,4}:){1,5}(:[0-9A-Fa-f]{1,4}){1,2})|(([0-9A-Fa-f]{1,4}:){1,4}(:[0-9A-Fa-f]{1,4}){1,3})|(([0-9A-Fa-f]{1,4}:){1,3}(:[0-9A-Fa-f]{1,4}){1,4})|(([0-9A-Fa-f]{1,4}:){1,2}(:[0-9A-Fa-f]{1,4}){1,5})|([0-9A-Fa-f]{1,4}(:([0-9A-Fa-f]{1,4}:){1,6}))|(:((:[0-9A-Fa-f]{1,4}){1,7}|:))|(([0-9A-Fa-f]{1,4}:){1,4}((25[0-5]|2[0-4][0-9]|[01]?[0-9]?[0-9])\.){3}(25[0-5]|2[0-4][0-9]|[01]?[0-9]?[0-9])))$'
IPV4_REGEX='^([0-9]{1,3}\.){3}[0-9]{1,3}$'

DYNU_HOSTNAME=$(bashio::config 'dynu_hostname' | tr '[:upper:]' '[:lower:]' | tr -d '[:space:]')
PASSWORD=$(bashio::config 'dynu_update_password' | xargs)
IP_DOMAINS=$(bashio::config 'ip_domains' | tr '[:upper:]' '[:lower:]' | xargs)
DYNU_API_KEY=$(bashio::config 'dynu_api_key' | xargs)

if bashio::config.true 'verbose'; then
    VERBOSE="true"
else
    VERBOSE="false"
fi
if bashio::config.true 'mqtt'; then
    USE_MQTT="true"
else
    USE_MQTT="false"
fi

MQTT_CONF="/run/dynu-mqtt.conf"

if [[ "$USE_MQTT" == "true" && -f "$MQTT_CONF" ]]; then
	MQTT_DATA=$(< "$MQTT_CONF")
	USE_MQTT=$(jq -r '.mqtt' <<< "$MQTT_DATA")
	MQTT_HOST=$(jq -r '.host' <<< "$MQTT_DATA")
	MQTT_PORT=$(jq -r '.port' <<< "$MQTT_DATA")
	MQTT_USER=$(jq -r '.user' <<< "$MQTT_DATA")
	MQTT_PASSWORD=$(jq -r '.password' <<< "$MQTT_DATA")
else
	USE_MQTT="false"
fi
	
mqtt_publish() {
    if [[ "$USE_MQTT" == "true" ]]; then
        local topic="$1"
        local payload="$2"

        if ! mosquitto_pub \
            -h "$MQTT_HOST" \
            -p "$MQTT_PORT" \
            -u "$MQTT_USER" \
            -P "$MQTT_PASSWORD" \
            -t "dynu_updater/$topic" \
            -m "$payload" \
            -r; then

            bashio::log.warning "Dynu DNS Updater - MQTT publish failed: Topic=$topic"
        fi
    fi
}

IP_CONF="/data/dynu-ip.json"

if [[ -f "$IP_CONF" ]]; then
	IP_DATA=$(< "$IP_CONF")
	CURRENT_IPV4=$(jq -r '.ipv4' <<< "$IP_DATA")
	CURRENT_IPV6=$(jq -r '.ipv6' <<< "$IP_DATA")
	LAST_IP_UPDATE=$(jq -r '.last' <<< "$IP_DATA")
else
	CURRENT_IPV4=""
	CURRENT_IPV6=""
	LAST_IP_UPDATE=""
fi

LAST_CHECK="true"
FIRST_RUN="true"
DYNU_IPV6_ENABLED="false"

load_api() {
	# Call Dynu API
	bashio::log.info "Dynu DNS Updater - Loading data from Dynu API"
	
	if DYNU_API_RESPONSE=$(curl -s --max-time 5 -H "API-Key: $DYNU_API_KEY" "https://api.dynu.com/v2/dns" 2>/dev/null); then
		HTTP_STATUS=$(
		jq -r --arg HOST "$DYNU_HOSTNAME" '
			.statusCode // empty
		' <<< "$DYNU_API_RESPONSE")

		case "$HTTP_STATUS" in
			200)
				bashio::log.info "Dynu DNS Updater - Dynu Domain API data retrieved."
				;;
			401)
				bashio::log.error "Dynu Updater Initialisation - Invalid Dynu API key"
				exit 1
				;;
			403)
				bashio::log.error "Dynu Updater Initialisation - Dynu Domain API key does not have permission to access DNS records"
				exit 1
				;;
			*)
				bashio::log.error "Dynu Updater Initialisation - Dynu Domain API returned HTTP $HTTP_STATUS"
				exit 1
				;;
		esac
	else
		bashio::log.error "Dynu Updater Initialisation - Dynu Domain API could not be loaded"
		exit 1
	fi
	
	# Retrieve Domain ID
	DYNU_DOMAIN_ID=$(
		jq -r --arg HOST "$DYNU_HOSTNAME" '
			.domains[]
			| select(.name == $HOST)
			| .id // empty
		' <<< "$DYNU_API_RESPONSE")

	if [[ -z "$DYNU_DOMAIN_ID" ]]; then
		bashio::log.error "Dynu DNS Updater - '$DYNU_HOSTNAME' not found in Dynu domain list"
		return 1
	fi
	
	# Retrieve IPv6 enabled status for domain
	DYNU_IPV6_ENABLED=$(jq -r --arg HOST "$DYNU_HOSTNAME" '
		.domains[]
		| select(.name == $HOST)
		| .ipv6 // "false"
	' <<<"$DYNU_API_RESPONSE")

	if DYNU_RECORDS=$(curl -s --max-time 5 -H "API-Key: $DYNU_API_KEY" "https://api.dynu.com/v2/dns/${DYNU_DOMAIN_ID}/record" 2>/dev/null); then
		HTTP_STATUS=$(
		jq -r --arg HOST "$DYNU_HOSTNAME" '
			.statusCode // empty
		' <<< "$DYNU_RECORDS")

		case "$HTTP_STATUS" in
			200)
				bashio::log.info "Dynu DNS Updater - Dynu Records API data retrieved."
				;;
			401)
				bashio::log.error "Dynu Updater Initialisation - Invalid Dynu API key"
				exit 1
				;;
			403)
				bashio::log.error "Dynu Updater Initialisation - Dynu API key does not have permission to access DNS records"
				exit 1
				;;
			*)
				bashio::log.error "Dynu Updater Initialisation - Dynu Records API returned HTTP $HTTP_STATUS"
				exit 1
				;;
		esac
	else
		bashio::log.error "Dynu Updater Initialisation - Dynu Records API could not be loaded"
		exit 1
	fi

	DYNU_HOSTNAMES=$(
	  jq -r '
		[.dnsRecords[]
		  | select(.recordType | IN("SOA","A","AAAA"))
		  | .hostname
		]
		| unique
		| sort_by( (split(".") | length), . )
		| .[]
	  ' <<< "$DYNU_RECORDS"
	)
	
	bashio::log.warning "Dynu Hosts: $DYNU_HOSTNAMES"

	DOMAIN_ARR=()
	for DOMAIN in $IP_DOMAINS; do
		if [[ "$DOMAIN" == "*.$DYNU_HOSTNAME" ]]; then
			read -r -a _tmp <<< "$IP_DOMAINS"
			if (( ${#_tmp[@]} != 1 )); then
				bashio::log.warning "Dynu DNS updater - Wildcard found, ignoring other DNS Update Hostnames"
			fi
			readarray -t DOMAIN_ARR <<< "$DYNU_HOSTNAMES"
			break
		elif ! echo "$DYNU_HOSTNAMES" | grep -qx "$DOMAIN"; then
			bashio::log.warning "Dynu DNS updater - $DOMAIN is not valid for this account"
			continue
		fi
		DOMAIN_ARR+=("$DOMAIN")
	done

	if [[ ${#DOMAIN_ARR[@]} -eq 0 ]]; then
		bashio::log.error "Dynu DNS updater - no valid domains found"
		return 1
	fi

	return 0
}

update_host() {
	bashio::log.info "Dynu DNS Updater - Updating $1"
				
	if [[ "$1" == "$DYNU_HOSTNAME" ]]; then
		COMMAND="https://api.dynu.com/nic/update?hostname=$DYNU_HOSTNAME&password=$PASSWORD&myip=$2"
	else
		COMMAND="https://api.dynu.com/nic/update?hostname=$DYNU_HOSTNAME&alias=${1%.$DYNU_HOSTNAME}&password=$PASSWORD&myip=$2"
	fi

	if [[ -n "$3" ]]; then
		COMMAND="${COMMAND}&myipv6=$3"
	fi
	
	if RESPONSE=$(curl -s --max-time 5 "$COMMAND" 2>/dev/null); then
		case "$RESPONSE" in
			good*|nochg*)
				bashio::log.info "Dynu DNS Updater - Success: $RESPONSE"
				return 0
				;;
			badauth)
				bashio::log.error "Dynu DNS Updater - Authentication failed (badauth)"
				return 1
				;;
			notfqdn)
				bashio::log.error "Dynu DNS Updater - Invalid hostname (notfqdn)"
				return 1
				;;
			nohost)
				bashio::log.error "Dynu DNS Updater - Host not found (nohost)"
				return 1
				;;
			numhost)
				bashio::log.error "Dynu DNS Updater - Too many hosts (numhost)"
				return 1
				;;
			abuse)
				bashio::log.error "Dynu DNS Updater - Account blocked (abuse)"
				return 1
				;;
			911)
				bashio::log.error "Dynu DNS Updater - Dynu server error (911)"
				return 1
				;;
			*)
				bashio::log.error "Dynu DNS Updater - Updating failed"
				return 1
				;;
		esac
	else
		bashio::log.error "Dynu DNS Updater - Updating failed"
		return 1
	fi
}

bashio::log.info "Dynu DNS Updater - Started"

while true; do
	if [[ "$FIRST_RUN" == "true" ]]; then
		if ! load_api; then
			bashio::log.error "Dynu DNS Updater - Unable to load Dynu API. Check log for details."
			exit 1
		fi
	fi
	
	if [[ "$VERBOSE" == "true" ]]; then
		bashio::log.info "Dynu DNS Updater - Checking current IP addresses."
	fi
	
	IPV4=""
	IPV6=""
	
	# --- Detect IPv4 ---
	if RAW=$(curl -4 -s --fail --max-time 5 https://checkip.dynu.com 2>/dev/null); then
		IPV4=$(echo "$RAW" | awk -F': ' '{print $2}' | tr -d '[:space:]')
	fi
	if [[ -z "$IPV4" || ! "$IPV4" =~ $IPV4_REGEX ]]; then
		bashio::log.warning "Dynu DNS Updater - dynu.com returned an invalid IPv4 address. Trying icanhazip.com"
		if RAW=$(curl -4 -s --fail --max-time 5 https://icanhazip.com 2>/dev/null); then
			IPV4=$(echo "$RAW" | tr -d '[:space:]')
		fi
		if [[ -z "$IPV4" || ! "$IPV4" =~ $IPV4_REGEX ]]; then
			bashio::log.warning "Dynu DNS Updater - icanhazip.com returned an invalid IPv4 address."
			IPV4=""
		fi
	fi
	
	# --- Detect IPv6 only if domain supports it ---
	if [[ "$DYNU_IPV6_ENABLED" == "true" ]]; then
		if RAW=$(curl -6 -s --fail --max-time 5 https://checkipv6.dynu.com 2>/dev/null); then
			IPV6=$(echo "$RAW" | awk -F': ' '{print $2}' | tr -d '[:space:]')
		fi
		if [[ -z "$IPV6" || ! "$IPV6" =~ $IPV6_REGEX ]]; then
			bashio::log.warning "Dynu DNS Updater - dynu.com returned an invalid IPv6 address. Trying icanhazip.com"
			if RAW=$(curl -6 -s --fail --max-time 5 https://icanhazip.com 2>/dev/null); then
				IPV6=$(echo "$RAW" | tr -d '[:space:]')
			fi
			if [[ -z "$IPV6" || ! "$IPV6" =~ $IPV6_REGEX ]]; then
				bashio::log.warning "Dynu DNS Updater - icanhazip.com returned an invalid IPv6 address."
				IPV6=""
			fi
		fi
	else
		IPV6=""
	fi
	
	if [[ "$VERBOSE" == "true" ]]; then
		bashio::log.info "Dynu DNS Updater - Current IPv4 address: $IPV4"
		if [[ "$DYNU_IPV6_ENABLED" == "true" ]]; then
			bashio::log.info "Dynu DNS Updater - Current IPv6 address: $IPV6"
		fi
	fi
	
	if [[ -n "$IPV4" ]] && [[ "$DYNU_IPV6_ENABLED" == "false" || -n "$IPV6" ]]; then
		mqtt_publish "wan_connection" "true"
		mqtt_publish "last_ip_check" "$(date -Is)"
		
		if [[ "$LAST_CHECK" == "false" || "$IPV4" != "$CURRENT_IPV4" || "$IPV6" != "$CURRENT_IPV6" ]]; then
			if [[ "$LAST_CHECK" == "false" ]]; then
				bashio::log.info "Dynu DNS Updater - Last update failed. Updating all hosts."
			else
				bashio::log.info "Dynu DNS Updater - IP address change detected. Updating all hosts"
			fi
						  
			UPDATED=${#DOMAIN_ARR[@]}
			
			for HOST in "${DOMAIN_ARR[@]}"; do
				if update_host "$HOST" "$IPV4" "$IPV6"; then
					UPDATED=$((UPDATED - 1))
				fi
			done
			
			if (( UPDATED > 0 )); then
				bashio::log.error "Dynu DNS Updater - $UPDATED could not be updated. Please check logs for more information"
				LAST_CHECK="false"
				mqtt_publish "update_status" "false"
			else
				CURRENT_IPV4="$IPV4"
				CURRENT_IPV6="$IPV6"
				LAST_IP_UPDATE="$(date -Is)"
				LAST_CHECK="true"
				
				mqtt_publish "update_status" "true"
				mqtt_publish "last_ip_update" "$(date -Is)"
				mqtt_publish "current_ipv4" "$IPV4"
				mqtt_publish "current_ipv6" "$IPV6"
				
				jq -n \
				  --arg js_ipv4 "$CURRENT_IPV4" \
				  --arg js_ipv6 "$CURRENT_IPV6" \
				  --arg js_last "$(date -Is)" \
				  '{
					  ipv4: $js_ipv4,
					  ipv6: $js_ipv6,
					  last: $js_last
				   }' \
				  > /data/dynu-ip.json
				  
				bashio::log.info "Dynu DNS Updater - All Hosts/Aliases were updated"
			fi
		elif [[ "$VERBOSE" == "true" || "$FIRST_RUN" == "true" ]]; then
			bashio::log.info "Dynu DNS Updater - No IP address change detected"
		fi

	else
		mqtt_publish "wan_connection" "false"
		bashio::log.error "Dynu DNS Updater - Unable to retrieve current IP address data (WAN Down?)"
	fi
	
	if [[ "$FIRST_RUN" == "true" ]]; then
		touch /run/dynu-dns-update.ready
		FIRST_RUN="false"
	fi
	
	if [[ "$VERBOSE" == "true" ]]; then
		bashio::log.info "Dynu DNS Updater - Checking and updating complete. Rechecking in 60 seconds"
	fi

	if [[ "LAST_CHECK" == "false" ]]; then
		sleep 300
	else
		sleep 60
	fi	
done