# dynu_updater/rootfs/etc/s6-overlay/s6-rc.d/finish/up
#!/usr/bin/with-contenv bash
# no-op
exit 0